# config/rootfiles/core/175/filelists/files
etc/rc.d/init.d/firewall
lib/firmware/amdgpu/dcn_3_1_4_dmcub.bin
lib/firmware/amdgpu/gc_11_0_1_imu.bin
lib/firmware/amdgpu/gc_11_0_1_me.bin
lib/firmware/amdgpu/gc_11_0_1_mec.bin
lib/firmware/amdgpu/gc_11_0_1_mes1.bin
lib/firmware/amdgpu/gc_11_0_1_pfp.bin
lib/firmware/amdgpu/gc_11_0_1_rlc.bin
lib/firmware/amdgpu/gc_11_0_4_imu.bin
lib/firmware/amdgpu/gc_11_0_4_me.bin
lib/firmware/amdgpu/gc_11_0_4_mec.bin
lib/firmware/amdgpu/gc_11_0_4_mes1.bin
lib/firmware/amdgpu/gc_11_0_4_mes.bin
lib/firmware/amdgpu/gc_11_0_4_pfp.bin
lib/firmware/amdgpu/gc_11_0_4_rlc.bin
lib/firmware/amdgpu/psp_13_0_11_ta.bin
lib/firmware/amdgpu/psp_13_0_11_toc.bin
lib/firmware/amdgpu/psp_13_0_4_ta.bin
lib/firmware/amdgpu/sdma_6_0_1.bin
lib/firmware/ath10k/QCA6174/hw3.0/firmware-sdio-6.bin
lib/firmware/ath10k/WCN3990/hw1.0/board-2.bin
lib/firmware/ath11k/IPQ5018/hw1.0/board-2.bin
lib/firmware/ath11k/IPQ5018/hw1.0/m3_fw.b00
lib/firmware/ath11k/IPQ5018/hw1.0/m3_fw.b01
lib/firmware/ath11k/IPQ5018/hw1.0/m3_fw.b02
lib/firmware/ath11k/IPQ5018/hw1.0/m3_fw.flist
lib/firmware/ath11k/IPQ5018/hw1.0/m3_fw.mdt
lib/firmware/ath11k/IPQ5018/hw1.0/q6_fw.b00
lib/firmware/ath11k/IPQ5018/hw1.0/q6_fw.b01
lib/firmware/ath11k/IPQ5018/hw1.0/q6_fw.b02
lib/firmware/ath11k/IPQ5018/hw1.0/q6_fw.b03
lib/firmware/ath11k/IPQ5018/hw1.0/q6_fw.b04
lib/firmware/ath11k/IPQ5018/hw1.0/q6_fw.b05
lib/firmware/ath11k/IPQ5018/hw1.0/q6_fw.b07
lib/firmware/ath11k/IPQ5018/hw1.0/q6_fw.b08
lib/firmware/ath11k/IPQ5018/hw1.0/q6_fw.b09
lib/firmware/ath11k/IPQ5018/hw1.0/q6_fw.b10
lib/firmware/ath11k/IPQ5018/hw1.0/q6_fw.b11
lib/firmware/ath11k/IPQ5018/hw1.0/q6_fw.b13
lib/firmware/ath11k/IPQ5018/hw1.0/q6_fw.b14
lib/firmware/ath11k/IPQ5018/hw1.0/q6_fw.flist
lib/firmware/ath11k/IPQ5018/hw1.0/q6_fw.mdt
lib/firmware/brcm/brcmfmac43241b4-sdio.Advantech-MICA-071.txt
lib/firmware/brcm/brcmfmac4356-pcie.Intel Corporation-CHERRYVIEW D1 PLATFORM.txt
lib/firmware/brcm/brcmfmac4356-pcie.Xiaomi Inc-Mipad2.txt
lib/firmware/i915/adlp_dmc.bin
lib/firmware/i915/mtl_dmc.bin
lib/firmware/iwlwifi-9260-th-b0-jf-b0-46.ucode
lib/firmware/iwlwifi-Qu-b0-hr-b0-77.ucode
lib/firmware/iwlwifi-Qu-b0-jf-b0-77.ucode
lib/firmware/iwlwifi-Qu-c0-jf-b0-77.ucode
lib/firmware/iwlwifi-QuZ-a0-hr-b0-77.ucode
lib/firmware/iwlwifi-so-a0-gf4-a0-72.ucode
lib/firmware/iwlwifi-so-a0-gf4-a0-74.ucode
lib/firmware/iwlwifi-so-a0-gf4-a0-81.ucode
lib/firmware/iwlwifi-so-a0-gf-a0-72.ucode
lib/firmware/iwlwifi-so-a0-gf-a0-74.ucode
lib/firmware/iwlwifi-so-a0-gf-a0-81.ucode
lib/firmware/iwlwifi-so-a0-gf-a0.pnvm
lib/firmware/iwlwifi-so-a0-hr-b0-81.ucode
lib/firmware/iwlwifi-ty-a0-gf-a0-72.ucode
lib/firmware/iwlwifi-ty-a0-gf-a0-74.ucode
lib/firmware/iwlwifi-ty-a0-gf-a0-81.ucode
lib/firmware/iwlwifi-ty-a0-gf-a0.pnvm
lib/firmware/mediatek/mt8192/scp.img
lib/firmware/mediatek/mt8195/scp.img
lib/firmware/mediatek/WIFI_MT7922_patch_mcu_1_1_hdr.bin
lib/firmware/mediatek/WIFI_MT7961_patch_mcu_1_2_hdr.bin
lib/firmware/mediatek/WIFI_RAM_CODE_MT7922_1.bin
lib/firmware/mediatek/WIFI_RAM_CODE_MT7961_1.bin
lib/firmware/nvidia/tu102/sec2/desc.bin
lib/firmware/nvidia/tu102/sec2/image.bin
lib/firmware/nvidia/tu102/sec2/sig.bin
lib/firmware/nvidia/tu116/sec2/desc.bin
lib/firmware/nvidia/tu116/sec2/image.bin
lib/firmware/nvidia/tu116/sec2/sig.bin
lib/firmware/qat_4xxx.bin
lib/firmware/qat_4xxx_mmp.bin
lib/firmware/rtlwifi/rtl8710bufw_SMIC.bin
lib/firmware/rtlwifi/rtl8710bufw_UMC.bin
lib/firmware/rtw88/rtw8822c_fw.bin
lib/firmware/rtw89/rtw8852b_fw-1.bin
usr/lib/firewall/rules.pl
