# config/rootfiles/oldcore/27/filelists/files
etc/rc.d/init.d/squid
etc/sysctl.conf
etc/squid/errors
usr/bin/squidclient
usr/lib/squid/auth
usr/lib/squid/cachemgr.cgi
usr/lib/squid/diskd-daemon
usr/lib/squid/errors/English
usr/lib/squid/errors/German
usr/lib/squid/icons
usr/lib/squid/logfile-daemon
usr/lib/squid/mib.txt
usr/lib/squid/msnt_auth
usr/lib/squid/ncsa_auth
usr/lib/squid/ntlm_auth
usr/lib/squid/smb_auth
usr/lib/squid/smb_auth.pl
usr/lib/squid/smb_auth.sh
usr/lib/squid/squid_ldap_auth
usr/lib/squid/unlinkd
usr/local/bin/squidctrl
usr/sbin/squid
var/ipfire/updatexlrator/bin/checkdeaddl
var/ipfire/updatexlrator/bin/checkup
var/ipfire/updatexlrator/bin/convert
var/ipfire/updatexlrator/bin/download
var/ipfire/updatexlrator/bin/lscache
var/ipfire/updatexlrator/updxlrator-lib.pl
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/cachemgr.cgi
srv/web/ipfire/cgi-bin/updatexlrator.cgi
var/ipfire/langs
lib/modules/KVER-ipfire/kernel/drivers/net/imq.ko
var/ipfire/qos/bin/makeqosscripts.pl
etc/rc.d/init.d/network
bin/ntfs-3g
lib/libntfs-3g.so
lib/libntfs-3g.so.49
lib/libntfs-3g.so.49.0.0
sbin/mount.ntfs-3g
usr/bin/ntfs-3g.probe
usr/lib/libntfs-3g.so
usr/lib/pkgconfig/libntfs-3g.pc
etc/madwifi.conf
etc/modprobe.d/madwifi
usr/bin/80211debug
usr/bin/80211stats
usr/bin/ath_info
usr/bin/athchans
usr/bin/athctrl
usr/bin/athdebug
usr/bin/athkey
usr/bin/athstats
usr/bin/madwifi-unload
usr/bin/wlanconfig
usr/bin/wpakey
usr/local/bin/rebuild-initrd
var/ipfire/menu.d/20-status.menu
srv/web/ipfire/cgi-bin/atm-status.cgi
usr/local/bin/backupiso
