# config/rootfiles/core/154/filelists/files
etc/system-release
etc/issue
etc/os-release
etc/rc.d/init.d/functions
etc/rc.d/init.d/sshd
etc/rc.d/init.d/unbound
etc/rc.d/init.d/wlanclient
lib/udev/network-hotplug-bridges
opt/pakfire/lib/functions.sh
srv/web/ipfire/cgi-bin/country.cgi
srv/web/ipfire/cgi-bin/credits.cgi
srv/web/ipfire/cgi-bin/ddns.cgi
srv/web/ipfire/cgi-bin/dhcp.cgi
srv/web/ipfire/cgi-bin/extrahd.cgi
srv/web/ipfire/cgi-bin/gpl.cgi
srv/web/ipfire/cgi-bin/hardwaregraphs.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/optionsfw.cgi
srv/web/ipfire/cgi-bin/pakfire.cgi
srv/web/ipfire/cgi-bin/remote.cgi
srv/web/ipfire/cgi-bin/services.cgi
srv/web/ipfire/cgi-bin/wirelessclient.cgi
usr/local/bin/ipsec-interfaces
usr/local/bin/sshctrl
usr/sbin/convert-snort
var/ipfire/extrahd/bin/extrahd.pl
var/ipfire/general-functions.pl
var/ipfire/ids-functions.pl
var/ipfire/network-functions.pl
var/ipfire/header.pl
var/ipfire/langs
