# config/rootfiles/core/185/filelists/files
etc/rc.d/init.d/functions
etc/sudoers.d/logwatch-mdadm
lib/firmware/brcm/BCM-0a5c-6410.hcd
lib/firmware/brcm/brcmfmac43012-sdio.bin
lib/firmware/brcm/brcmfmac43012-sdio.clm_blob
lib/firmware/brcm/brcmfmac43430-sdio.clm_blob
lib/firmware/brcm/brcmfmac43430-sdio.raspberrypi,model-zero-w.txt
lib/firmware/brcm/brcmfmac43430-sdio.sinovoip,bpi-m2-plus.txt
lib/firmware/brcm/brcmfmac43430-sdio.sinovoip,bpi-m2-ultra.txt
lib/firmware/brcm/brcmfmac43430-sdio.sinovoip,bpi-m2-zero.txt
lib/firmware/brcm/brcmfmac43430-sdio.sinovoip,bpi-m3.txt
lib/firmware/brcm/brcmfmac43455-sdio.clm_blob
lib/firmware/brcm/brcmfmac43455-sdio.raspberrypi,3-model-a-plus.txt
lib/firmware/brcm/brcmfmac43455-sdio.Raspberry Pi Foundation-Raspberry Pi 4 Model B.txt
lib/firmware/brcm/brcmfmac43455-sdio.Raspberry Pi Foundation-Raspberry Pi Compute Module 4.txt
lib/firmware/brcm/brcmfmac4354-sdio.clm_blob
lib/firmware/brcm/brcmfmac4356-pcie.clm_blob
lib/firmware/brcm/brcmfmac4356-sdio.clm_blob
lib/firmware/brcm/brcmfmac4356-sdio.khadas,vim2.txt
lib/firmware/brcm/brcmfmac43570-pcie.clm_blob
lib/firmware/brcm/brcmfmac4373-sdio.clm_blob
lib/firmware/brcm/brcmfmac54591-pcie.bin
lib/firmware/brcm/brcmfmac54591-pcie.clm_blob
lib/firmware/cxgb4/t4-config.txt
lib/firmware/cxgb4/t5-config.txt
lib/firmware/cxgb4/t6-config.txt
lib/firmware/intel/ice/ddp/ice.pkg
lib/firmware/netronome/flower/nic_AMDA0058-0011_1x100.nffw
lib/firmware/netronome/flower/nic_AMDA0058-0011_2x40.nffw
lib/firmware/netronome/flower/nic_AMDA0058-0011_4x10_1x40.nffw
lib/firmware/netronome/flower/nic_AMDA0058-0011_8x10.nffw
lib/firmware/netronome/flower/nic_AMDA0058-0012_1x100.nffw
lib/firmware/netronome/flower/nic_AMDA0058-0012_2x40.nffw
lib/firmware/netronome/flower/nic_AMDA0058-0012_4x10_1x40.nffw
lib/firmware/netronome/flower/nic_AMDA0058-0012_8x10.nffw
lib/firmware/netronome/flower/nic_AMDA0078-0011_1x100.nffw
lib/firmware/netronome/flower/nic_AMDA0078-0011_2x40.nffw
lib/firmware/netronome/flower/nic_AMDA0078-0011_4x10_1x40.nffw
lib/firmware/netronome/flower/nic_AMDA0078-0011_8x10.nffw
lib/firmware/netronome/flower/nic_AMDA0078-0012_1x100.nffw
lib/firmware/netronome/flower/nic_AMDA0078-0012_2x40.nffw
lib/firmware/netronome/flower/nic_AMDA0078-0012_4x10_1x40.nffw
lib/firmware/netronome/flower/nic_AMDA0078-0012_8x10.nffw
lib/firmware/nvidia/tegra124/vic.bin
lib/firmware/nvidia/tegra186/vic.bin
lib/firmware/nvidia/tegra210/vic.bin
srv/web/ipfire/cgi-bin/dhcp.cgi
srv/web/ipfire/cgi-bin/dns.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/time.cgi
var/ipfire/backup/bin/backup.pl
var/ipfire/main/manualpages
var/ipfire/ovpn/openssl/ovpn.cnf
