# config/rootfiles/core/197/filelists/files
etc/httpd/conf/vhosts.d/ipfire-interface-ssl.conf
etc/rc.d/init.d/cpupower
etc/rc.d/init.d/firewall
etc/rc.d/init.d/functions
etc/rc.d/init.d/openvpn-n2n
etc/rc.d/init.d/openvpn-rw
etc/rc.d/rc0.d/K10openvpn-rw
etc/rc.d/rc0.d/K11openvpn-n2n
etc/rc.d/rc3.d/S50openvpn-rw
etc/rc.d/rc3.d/S51openvpn-n2n
etc/rc.d/rc6.d/K10openvpn-rw
etc/rc.d/rc6.d/K11openvpn-n2n
etc/rc.d/rcsysinit.d/S46cpupower
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/services.cgi
srv/web/ipfire/cgi-bin/vulnerabilities.cgi
srv/web/ipfire/html/themes/ipfire/include/css/style.css
usr/local/bin/openvpnctrl
var/ipfire/backup/bin/backup.pl
var/ipfire/general-functions.pl
var/ipfire/header.pl
var/ipfire/langs/list
var/ipfire/menu.d/20-status.menu
var/ipfire/menu.d/40-services.menu
var/ipfire/main/manualpages
