# config/rootfiles/oldcore/173/filelists/files
etc/ppp/ip-up
etc/rc.d/init.d/networking/dhcpcd.exe
etc/rc.d/init.d/networking/functions.network
etc/rc.d/init.d/networking/red
lib/firmware/amdgpu/aldebaran_mec2.bin
lib/firmware/amdgpu/aldebaran_mec.bin
lib/firmware/amdgpu/aldebaran_rlc.bin
lib/firmware/amdgpu/aldebaran_sjt_mec2.bin
lib/firmware/amdgpu/aldebaran_sjt_mec.bin
lib/firmware/amdgpu/aldebaran_smc.bin
lib/firmware/amdgpu/arcturus_mec2.bin
lib/firmware/amdgpu/arcturus_mec.bin
lib/firmware/amdgpu/arcturus_ta.bin
lib/firmware/amdgpu/beige_goby_ce.bin
lib/firmware/amdgpu/beige_goby_dmcub.bin
lib/firmware/amdgpu/beige_goby_me.bin
lib/firmware/amdgpu/beige_goby_mec2.bin
lib/firmware/amdgpu/beige_goby_mec.bin
lib/firmware/amdgpu/beige_goby_pfp.bin
lib/firmware/amdgpu/beige_goby_smc.bin
lib/firmware/amdgpu/beige_goby_sos.bin
lib/firmware/amdgpu/beige_goby_ta.bin
lib/firmware/amdgpu/dcn_3_2_0_dmcub.bin
lib/firmware/amdgpu/dimgrey_cavefish_dmcub.bin
lib/firmware/amdgpu/dimgrey_cavefish_smc.bin
lib/firmware/amdgpu/dimgrey_cavefish_sos.bin
lib/firmware/amdgpu/dimgrey_cavefish_ta.bin
lib/firmware/amdgpu/dimgrey_cavefish_vcn.bin
lib/firmware/amdgpu/gc_10_3_7_rlc.bin
lib/firmware/amdgpu/gc_11_0_0_imu.bin
lib/firmware/amdgpu/gc_11_0_0_me.bin
lib/firmware/amdgpu/gc_11_0_0_mec.bin
lib/firmware/amdgpu/gc_11_0_0_mes1.bin
lib/firmware/amdgpu/gc_11_0_0_mes.bin
lib/firmware/amdgpu/gc_11_0_0_pfp.bin
lib/firmware/amdgpu/gc_11_0_0_rlc.bin
lib/firmware/amdgpu/green_sardine_asd.bin
lib/firmware/amdgpu/green_sardine_dmcub.bin
lib/firmware/amdgpu/green_sardine_ta.bin
lib/firmware/amdgpu/green_sardine_vcn.bin
lib/firmware/amdgpu/navi10_asd.bin
lib/firmware/amdgpu/navi10_ta.bin
lib/firmware/amdgpu/navi10_vcn.bin
lib/firmware/amdgpu/navi12_asd.bin
lib/firmware/amdgpu/navi12_ce.bin
lib/firmware/amdgpu/navi12_me.bin
lib/firmware/amdgpu/navi12_mec2.bin
lib/firmware/amdgpu/navi12_mec.bin
lib/firmware/amdgpu/navi12_pfp.bin
lib/firmware/amdgpu/navi12_sdma1.bin
lib/firmware/amdgpu/navi12_sdma.bin
lib/firmware/amdgpu/navi12_sos.bin
lib/firmware/amdgpu/navi12_ta.bin
lib/firmware/amdgpu/navi12_vcn.bin
lib/firmware/amdgpu/navi14_asd.bin
lib/firmware/amdgpu/navi14_ta.bin
lib/firmware/amdgpu/navi14_vcn.bin
lib/firmware/amdgpu/navy_flounder_dmcub.bin
lib/firmware/amdgpu/navy_flounder_smc.bin
lib/firmware/amdgpu/navy_flounder_sos.bin
lib/firmware/amdgpu/navy_flounder_ta.bin
lib/firmware/amdgpu/navy_flounder_vcn.bin
lib/firmware/amdgpu/picasso_asd.bin
lib/firmware/amdgpu/picasso_ce.bin
lib/firmware/amdgpu/picasso_me.bin
lib/firmware/amdgpu/picasso_mec2.bin
lib/firmware/amdgpu/picasso_mec.bin
lib/firmware/amdgpu/picasso_pfp.bin
lib/firmware/amdgpu/picasso_ta.bin
lib/firmware/amdgpu/picasso_vcn.bin
lib/firmware/amdgpu/psp_13_0_0_sos.bin
lib/firmware/amdgpu/psp_13_0_0_ta.bin
lib/firmware/amdgpu/psp_13_0_5_asd.bin
lib/firmware/amdgpu/psp_13_0_5_ta.bin
lib/firmware/amdgpu/psp_13_0_8_ta.bin
lib/firmware/amdgpu/raven2_vcn.bin
lib/firmware/amdgpu/raven_asd.bin
lib/firmware/amdgpu/raven_vcn.bin
lib/firmware/amdgpu/renoir_asd.bin
lib/firmware/amdgpu/renoir_ce.bin
lib/firmware/amdgpu/renoir_dmcub.bin
lib/firmware/amdgpu/renoir_me.bin
lib/firmware/amdgpu/renoir_mec2.bin
lib/firmware/amdgpu/renoir_mec.bin
lib/firmware/amdgpu/renoir_pfp.bin
lib/firmware/amdgpu/renoir_ta.bin
lib/firmware/amdgpu/renoir_vcn.bin
lib/firmware/amdgpu/sdma_5_2_7.bin
lib/firmware/amdgpu/sdma_6_0_0.bin
lib/firmware/amdgpu/sienna_cichlid_ce.bin
lib/firmware/amdgpu/sienna_cichlid_dmcub.bin
lib/firmware/amdgpu/sienna_cichlid_me.bin
lib/firmware/amdgpu/sienna_cichlid_mec2.bin
lib/firmware/amdgpu/sienna_cichlid_mec.bin
lib/firmware/amdgpu/sienna_cichlid_pfp.bin
lib/firmware/amdgpu/sienna_cichlid_smc.bin
lib/firmware/amdgpu/sienna_cichlid_sos.bin
lib/firmware/amdgpu/sienna_cichlid_ta.bin
lib/firmware/amdgpu/sienna_cichlid_vcn.bin
lib/firmware/amdgpu/smu_13_0_0.bin
lib/firmware/amdgpu/vangogh_asd.bin
lib/firmware/amdgpu/vangogh_ce.bin
lib/firmware/amdgpu/vangogh_dmcub.bin
lib/firmware/amdgpu/vangogh_me.bin
lib/firmware/amdgpu/vangogh_mec2.bin
lib/firmware/amdgpu/vangogh_mec.bin
lib/firmware/amdgpu/vangogh_pfp.bin
lib/firmware/amdgpu/vcn_4_0_0.bin
lib/firmware/amdgpu/vega10_asd.bin
lib/firmware/amdgpu/vega10_ce.bin
lib/firmware/amdgpu/vega10_me.bin
lib/firmware/amdgpu/vega10_mec2.bin
lib/firmware/amdgpu/vega10_mec.bin
lib/firmware/amdgpu/vega10_pfp.bin
lib/firmware/amdgpu/vega20_asd.bin
lib/firmware/amdgpu/vega20_ce.bin
lib/firmware/amdgpu/vega20_me.bin
lib/firmware/amdgpu/vega20_mec2.bin
lib/firmware/amdgpu/vega20_mec.bin
lib/firmware/amdgpu/vega20_pfp.bin
lib/firmware/amdgpu/yellow_carp_ce.bin
lib/firmware/amdgpu/yellow_carp_dmcub.bin
lib/firmware/amdgpu/yellow_carp_me.bin
lib/firmware/amdgpu/yellow_carp_mec2.bin
lib/firmware/amdgpu/yellow_carp_mec.bin
lib/firmware/amdgpu/yellow_carp_pfp.bin
lib/firmware/amdgpu/yellow_carp_ta.bin
lib/firmware/brcm/brcmfmac43241b0-sdio.bin
lib/firmware/brcm/brcmfmac4329-sdio.bin
lib/firmware/brcm/brcmfmac4330-sdio.bin
lib/firmware/brcm/brcmfmac4358-pcie.bin
lib/firmware/i915/dg2_dmc_ver2_08.bin
lib/firmware/i915/mtl_dmc_ver2_10.bin
lib/firmware/mediatek/mt7915_rom_patch.bin
lib/firmware/mediatek/mt7915_wa.bin
lib/firmware/mediatek/mt7915_wm.bin
lib/firmware/mediatek/mt7916_rom_patch.bin
lib/firmware/mediatek/mt7916_wa.bin
lib/firmware/mediatek/mt7916_wm.bin
lib/firmware/mediatek/mt7986_rom_patch.bin
lib/firmware/mediatek/mt7986_rom_patch_mt7975.bin
lib/firmware/mediatek/mt7986_wa.bin
lib/firmware/mediatek/mt7986_wm.bin
lib/firmware/mediatek/mt7986_wm_mt7975.bin
lib/firmware/mediatek/mt7986_wo_0.bin
lib/firmware/mediatek/mt7986_wo_1.bin
lib/firmware/mediatek/WIFI_MT7961_patch_mcu_1_2_hdr.bin
lib/firmware/mediatek/WIFI_RAM_CODE_MT7961_1.bin
lib/firmware/qca/nvm_usb_00190200.bin
lib/firmware/qca/rampatch_usb_00190200.bin
lib/udev/network-hotplug-rename
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/pppsetup.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/services.cgi
usr/lib/firewall/rules.pl
var/ipfire/dhcpc/dhcpcd.conf
var/ipfire/general-functions.pl
