# config/rootfiles/core/120/filelists/files
etc/system-release
etc/issue
etc/sysctl.conf
etc/fcron.daily/openvpn-crl-updater
etc/rc.d/init.d/dhcp
srv/web/ipfire/cgi-bin/ovpnmain.cgi
usr/lib/python2.7/lib-dynload/_hashlib.so
usr/lib/python2.7/lib-dynload/_ssl.so
var/ipfire/langs
var/ipfire/ovpn/openssl/ovpn.cnf
