# config/rootfiles/core/120/filelists/files
etc/system-release
etc/issue
etc/sysctl.conf
etc/fcron.daily/openvpn-crl-updater
etc/rc.d/init.d/dhcp
etc/rc.d/init.d/pakfire
etc/rc.d/init.d/unbound
etc/rc.d/rcsysinit.d/S81pakfire
opt/pakfire/lib/functions.pl
opt/pakfire/pakfire
opt/pakfire/pakfire-2007.key
opt/pakfire/pakfire-2018.key
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/qos.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
usr/lib/python2.7/lib-dynload/_hashlib.so
usr/lib/python2.7/lib-dynload/_ssl.so
var/ipfire/langs
var/ipfire/ovpn/openssl/ovpn.cnf
var/ipfire/qos/bin/makeqosscripts.pl
