# config/rootfiles/core/120/filelists/files
etc/system-release
etc/issue
etc/fcron.daily/openvpn-crl-updater
etc/rc.d/init.d/dhcp
srv/web/ipfire/cgi-bin/ovpnmain.cgi
var/ipfire/langs
var/ipfire/ovpn/openssl/ovpn.cnf
