# config/rootfiles/core/164/filelists/files
etc/collectd.conf
etc/rc.d/init.d/firewall
etc/rc.d/init.d/squid
etc/suricata/suricata.yaml
etc/sysctl.conf
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/optionsfw.cgi
srv/web/ipfire/cgi-bin/pakfire.cgi
srv/web/ipfire/html/include/pakfire.js
usr/sbin/convert-ids-multiple-providers
usr/sbin/convert-snort
var/ipfire/backup/bin/backup.pl
var/ipfire/backup/include
var/ipfire/graphs.pl
var/ipfire/ids-functions.pl
var/ipfire/urlfilter/autoupdate/autoupdate.urls
