# config/rootfiles/core/164/filelists/files
etc/collectd.conf
etc/rc.d/init.d/squid
var/ipfire/urlfilter/autoupdate/autoupdate.urls
