# config/rootfiles/core/125/filelists/files
etc/system-release
etc/issue
etc/ssh/sshd_config
etc/sysctl.conf
etc/rc.d/init.d/sysklogd
etc/rc.d/init.d/udev
srv/web/ipfire/cgi-bin/credits.cgi
srv/web/ipfire/cgi-bin/hardwaregraphs.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/media.cgi
srv/web/ipfire/cgi-bin/pakfire.cgi
usr/sbin/firewall-policy
usr/sbin/updxlrator
usr/local/bin/backupiso
usr/local/bin/makegraphs
var/ipfire/backup/bin/backup.pl
var/ipfire/backup/include
var/ipfire/langs
