# config/rootfiles/core/162/filelists/files
etc/rc.d/init.d/firewall
etc/rc.d/init.d/mountkernfs
etc/rc.d/init.d/unbound
srv/web/ipfire/cgi-bin/connections.cgi
srv/web/ipfire/cgi-bin/fwhosts.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/logs.cgi/log.dat
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
srv/web/ipfire/html/themes/ipfire/include/css/style.css
srv/web/ipfire/html/themes/ipfire/include/functions.pl
var/ipfire/general-functions.pl
var/ipfire/header.pl
var/ipfire/location-functions.pl
var/ipfire/main/manualpages
var/ipfire/qos/bin/makeqosscripts.pl
