# config/rootfiles/oldcore/153/filelists/files
etc/issue
etc/ntp.conf
etc/os-release
etc/rc.d/init.d/network
etc/rc.d/init.d/unbound
etc/system-release
etc/sysctl.conf
lib/udev/network-hotplug-bridges
srv/web/ipfire/cgi-bin/connections.cgi
srv/web/ipfire/cgi-bin/country.cgi
srv/web/ipfire/cgi-bin/credits.cgi
srv/web/ipfire/cgi-bin/dns.cgi
srv/web/ipfire/cgi-bin/ddns.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/ipinfo.cgi
srv/web/ipfire/cgi-bin/location-block.cgi
srv/web/ipfire/cgi-bin/logs.cgi/firewalllog.dat
srv/web/ipfire/cgi-bin/logs.cgi/firewalllogcountry.dat
srv/web/ipfire/cgi-bin/logs.cgi/firewalllogip.dat
srv/web/ipfire/cgi-bin/logs.cgi/showrequestfromcountry.dat
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/remote.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
srv/web/ipfire/cgi-bin/zoneconf.cgi
srv/web/ipfire/html/include/zoneconf.js
usr/lib/firewall/rules.pl
var/ipfire/location-functions.pl
var/ipfire/langs
