# config/rootfiles/core/13/files
etc/init.d/mISDN
etc/modprobe.d/blacklist
usr/bin/bc
usr/bin/dc
opt/pakfire/lib/functions.pl
var/ipfire/urlfilter/bin/prebuild.pl
srv/web/ipfire/cgi-bin/urlfilter.cgi
srv/web/ipfire/cgi-bin/updatexlrator.cgi
srv/web/ipfire/html/images/updxl-src-avg.gif
usr/sbin/updxlrator
