# config/rootfiles/core/13/files
etc/init.d/mISDN
etc/modprobe.d/blacklist
usr/bin/bc
usr/bin/dc
opt/pakfire/lib/functions.pl
