# config/rootfiles/core/61/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/collectd
opt/pakfire/lib/functions.sh
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
var/ipfire/general-functions.pl
var/ipfire/langs
usr/sbin/redirect_wrapper
