# config/rootfiles/oldcore/198/filelists/files
etc/rc.d/init.d/suricata
var/ipfire/graphs.pl
var/ipfire/ids-functions.pl
