# config/rootfiles/core/5/files
lib/modules/2.6.16.57-ipfire/kernel/fs/nfsd/nfsd.ko
lib/modules/2.6.16.57-ipfire-smp/kernel/fs/nfsd/nfsd.ko
usr/local/bin/qosctrl
etc/init.d/squid
var/ipfire/langs
var/ipfire/graphs.pl
