# config/rootfiles/core/111/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/networking/dhcpcd.exe
etc/rc.d/init.d/networking/red
etc/rc.d/init.d/wlanclient
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/ipinfo.cgi
srv/web/ipfire/cgi-bin/iptables.cgi
srv/web/ipfire/cgi-bin/logs.cgi/firewalllogcountry.dat
srv/web/ipfire/cgi-bin/logs.cgi/firewalllogip.dat
srv/web/ipfire/cgi-bin/logs.cgi/firewalllogport.dat
srv/web/ipfire/cgi-bin/logs.cgi/showrequestfromcountry.dat
srv/web/ipfire/cgi-bin/logs.cgi/showrequestfromip.dat
srv/web/ipfire/cgi-bin/logs.cgi/showrequestfromport.dat
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
srv/web/ipfire/cgi-bin/wirelessclient.cgi
srv/web/ipfire/html/images/back.png
usr/local/bin/getipstat
usr/local/bin/ipsecctrl
var/ipfire/langs
var/ipfire/graphs.pl
var/ipfire/network-functions.pl
var/ipfire/qos/bin/makeqosscripts.pl
