# config/rootfiles/core/110/filelists/files
etc/system-release
etc/issue
etc/httpd/conf/server-tuning.conf
etc/rc.d/init.d/unbound
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
usr/lib/libssp.so.0
usr/lib/libssp.so.0.0.0
usr/local/bin/xt_geoip_update
var/ipfire/langs
