# config/rootfiles/core/3/files
etc/init.d/networking/red
opt/pakfire/lib/functions.pl
opt/pakfire/pakfire
srv/web/ipfire/cgi-bin/pakfire.cgi
srv/web/ipfire/cgi-bin/updatexlrator.cgi
srv/web/ipfire/cgi-bin/urlfilter.cgi
srv/web/ipfire/html/images/updxl-src-apple.gif
usr/lib/ipsec/_realsetup
usr/sbin/dnsmasq
usr/sbin/updxlrator
var/ipfire/langs
var/ipfire/graphs.pl
var/spool/cron/root.orig
