# config/rootfiles/core/183/filelists/files
etc/rc.d/init.d/mountfs
etc/rc.d/init.d/networking/red.up/60-collectd
srv/web/ipfire/cgi-bin/dhcp.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/logs.cgi/firewalllog.dat
usr/local/bin/backupiso
