# config/rootfiles/core/190/filelists/files
etc/collectd.conf
etc/collectd.d
etc/httpd/conf/vhosts.d/ipfire-interface-ssl.conf
etc/rc.d/helper/aws-setup
etc/rc.d/helper/azure-setup
etc/rc.d/helper/exoscale-setup
etc/rc.d/helper/gcp-setup
etc/rc.d/helper/oci-setup
etc/rc.d/init.d/apache
etc/rc.d/init.d/cleanfs
etc/rc.d/init.d/collectd
etc/rc.d/init.d/firewall
etc/rc.d/init.d/networking/functions.network
etc/rc.d/init.d/squid
etc/rc.d/init.d/suricata
lib/udev/network-hotplug-bridges
srv/web/ipfire/cgi-bin/getrrdimage.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/logs.cgi/log.dat
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
usr/bin/suricata-watcher
usr/lib/firewall/rules.pl
usr/lib/perl5/5.36.0/xxxMACHINExxx-linux-thread-multi/Compress/Raw/Zlib.pm
usr/local/bin/openvpnctrl
usr/local/bin/sshctrl
var/ipfire/backup/include
var/ipfire/graphs.pl
var/ipfire/network-functions.pl
