# config/rootfiles/oldcore/122/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/acpid
etc/rc.d/init.d/collectd
etc/rc.d/init.d/firstsetup
etc/rc.d/init.d/leds
etc/rc.d/init.d/partresize
etc/rc.d/init.d/random
etc/rc.d/rc0.d/K87acpid
etc/rc.d/rc3.d/S12acpid
etc/rc.d/rc6.d/K87acpid
srv/web/ipfire/cgi-bin/entropy.cgi
srv/web/ipfire/cgi-bin/firewall.cgi
srv/web/ipfire/cgi-bin/fwhosts.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/remote.cgi
usr/lib/firewall/firewall-lib.pl
var/ipfire/langs
var/ipfire/qos/bin/makeqosscripts.pl
