# config/rootfiles/core/188/filelists/files
etc/rc.d/init.d/connectd
etc/rc.d/init.d/functions
etc/rc.d/init.d/networking/functions.network
etc/rc.d/init.d/networking/red
etc/rc.d/init.d/suricata
lib/udev/network-aqm
lib/udev/network-hotplug-vlan
srv/web/ipfire/cgi-bin/dhcp.cgi
srv/web/ipfire/cgi-bin/services.cgi
srv/web/ipfire/cgi-bin/system.cgi
srv/web/ipfire/cgi-bin/vulnerabilities.cgi
var/ipfire/backup/bin/backup.pl
var/ipfire/general-functions.pl
var/ipfire/graphs.pl
var/ipfire/header.pl
var/ipfire/lang.pl
var/ipfire/network-functions.pl
