# config/rootfiles/core/6/files
lib/modules/2.6.16.57-ipfire/kernel/fs/fuse
lib/modules/2.6.16.57-ipfire-smp/kernel/fs/fuse
lib/modules/2.6.16.57-ipfire/kernel/drivers/net/r8169.ko
lib/modules/2.6.16.57-ipfire/kernel/drivers/net/r8168.ko
lib/modules/2.6.16.57-ipfire/kernel/drivers/usb/net/mcs7830.ko
lib/modules/2.6.16.57-ipfire-smp/kernel/drivers/ieee1394/sbp2.ko
lib/modules/2.6.16.57-ipfire-smp/kernel/drivers/net/r8169.ko
lib/modules/2.6.16.57-ipfire-smp/kernel/drivers/net/r8168.ko
lib/modules/2.6.16.57-ipfire-smp/kernel/drivers/usb/net/mcs7830.ko
lib/modules/2.6.16.57-ipfire-smp/kernel/drivers/ieee1394/sbp2.ko
/sbin/kudzu
var/ipfire/langs
srv/web/ipfire/cgi-bin/outgoingfw.cgi
etc/init.d/squid
#bin/ntfs-3g
#lib/libntfs-3g.so
#lib/libntfs-3g.so.21
#lib/libntfs-3g.so.21.0.0
#sbin/mount.ntfs-3g
