# config/sarg/cron.hourly
#!/bin/bash

/usr/sbin/update-sarg-reports today >/dev/null 2>&1

exit 0
