# config/sarg/cron.monthly
#!/bin/bash

/usr/sbin/update-sarg-reports monthly >/dev/null 2>&1

exit 0
