# config/rootfiles/core/115/filelists/files
etc/system-release
etc/issue
etc/httpd/conf/vhosts.d/captive.conf
etc/httpd/conf/vhosts.d/ipfire-interface.conf
etc/httpd/conf/vhosts.d/ipfire-interface-ssl.conf
etc/rc.d/init.d/apache
etc/rc.d/init.d/firewall
srv/web/ipfire/cgi-bin/captive/index.cgi
srv/web/ipfire/cgi-bin/captive/logo.cgi
srv/web/ipfire/cgi-bin/captive/redirect.cgi
srv/web/ipfire/cgi-bin/captive.cgi
srv/web/ipfire/cgi-bin/extrahd.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
srv/web/ipfire/html/captive
usr/bin/captive-cleanup
usr/local/bin/backupiso
usr/local/bin/captivectrl
usr/local/bin/wirelessctrl
var/ipfire/backup/include
var/ipfire/captive
var/ipfire/extrahd/bin/extrahd.pl
var/ipfire/general-functions.pl
var/ipfire/lang.pl
var/ipfire/langs
var/ipfire/menu.d/30-network.menu
var/ipfire/modem-lib.pl
var/ipfire/network-functions.pl
var/ipfire/ovpn/openssl/ovpn.cnf
var/log/httpd/captive
