# config/rootfiles/core/195/filelists/files
etc/fcron.cyclic/wg-dynamic
etc/rc.d/init.d/firewall
etc/rc.d/init.d/networking/functions.network
etc/rc.d/init.d/wireguard
lib/udev/network-aqm
opt/pakfire/lib/functions.pl
srv/web/ipfire/cgi-bin/ddns.cgi
srv/web/ipfire/cgi-bin/firewall.cgi
srv/web/ipfire/cgi-bin/fwhosts.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/pakfire.cgi
srv/web/ipfire/cgi-bin/services.cgi
srv/web/ipfire/cgi-bin/wireguard.cgi
srv/web/ipfire/html/themes/ipfire/include/css/style.css
usr/lib/firewall/firewall-lib.pl
usr/local/bin/wireguardctrl
usr/sbin/firewall-policy
var/ipfire/backup/bin/backup.pl
var/ipfire/general-functions.pl
var/ipfire/header.pl
var/ipfire/http-client-functions.pl
var/ipfire/ids-functions.pl
var/ipfire/ipblocklist/sources
var/ipfire/menu.d/40-services.menu
var/ipfire/wireguard-functions.pl
