# config/rootfiles/core/152/filelists/files
etc/system-release
etc/issue
etc/rc.d/helpers/aws-setup
etc/rc.d/helpers/azure-setup
etc/rc.d/helpers/exoscale-setup
etc/rc.d/helpers/gcp-setup
etc/rc.d/helpers/oci-setup
etc/rc.d/init.d/networking/any
etc/rc.d/init.d/networking/red
srv/web/ipfire/cgi-bin/credits.cgi
usr/local/bin/setaliases
