# config/rootfiles/core/180/filelists/files
etc/logrotate.conf
etc/rc.d/init.d/network
etc/rc.d/init.d/networking/any
etc/rc.d/init.d/networking/functions.network
etc/rc.d/init.d/networking/red
etc/rc.d/init.d/networking/red.up/99-pakfire-update
lib/udev/rules.d/61-extrahd.rules
opt/pakfire/lib/functions.pl
srv/web/ipfire/cgi-bin/extrahd.cgi
srv/web/ipfire/cgi-bin/index.cgi
var/ipfire/main/manualpages
var/ipfire/extrahd/bin/extrahd.pl
