# config/rootfiles/core/70/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/squid
usr/local/bin/openvpnctrl
etc/init.d/firewall
