# config/rootfiles/core/70/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/firewall
etc/rc.d/init.d/network-vlans
etc/rc.d/init.d/squid
usr/local/bin/openvpnctrl
srv/web/ipfire/cgi-bin/ddns.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
usr/lib/squid/errors/tr/ERR_ACCESS_DENIED
usr/lib/squid/errors/tr/ERR_ACL_TIME_QUOTA_EXCEEDED
usr/lib/squid/errors/tr/ERR_AGENT_CONFIGURE
usr/lib/squid/errors/tr/ERR_AGENT_WPAD
usr/lib/squid/errors/tr/ERR_CACHE_ACCESS_DENIED
usr/lib/squid/errors/tr/ERR_CACHE_MGR_ACCESS_DENIED
usr/lib/squid/errors/tr/ERR_CANNOT_FORWARD
usr/lib/squid/errors/tr/ERR_CONNECT_FAIL
usr/lib/squid/errors/tr/ERR_DIR_LISTING
usr/lib/squid/errors/tr/ERR_DNS_FAIL
usr/lib/squid/errors/tr/ERR_ESI
usr/lib/squid/errors/tr/ERR_FORWARDING_DENIED
usr/lib/squid/errors/tr/ERR_FTP_DISABLED
usr/lib/squid/errors/tr/ERR_FTP_FAILURE
usr/lib/squid/errors/tr/ERR_FTP_FORBIDDEN
usr/lib/squid/errors/tr/ERR_FTP_NOT_FOUND
usr/lib/squid/errors/tr/ERR_FTP_PUT_CREATED
usr/lib/squid/errors/tr/ERR_FTP_PUT_ERROR
usr/lib/squid/errors/tr/ERR_FTP_PUT_MODIFIED
usr/lib/squid/errors/tr/ERR_FTP_UNAVAILABLE
usr/lib/squid/errors/tr/ERR_GATEWAY_FAILURE
usr/lib/squid/errors/tr/ERR_ICAP_FAILURE
usr/lib/squid/errors/tr/ERR_INVALID_REQ
usr/lib/squid/errors/tr/ERR_INVALID_RESP
usr/lib/squid/errors/tr/ERR_INVALID_URL
usr/lib/squid/errors/tr/ERR_LIFETIME_EXP
usr/lib/squid/errors/tr/ERR_NO_RELAY
usr/lib/squid/errors/tr/ERR_ONLY_IF_CACHED_MISS
usr/lib/squid/errors/tr/ERR_PRECONDITION_FAILED
usr/lib/squid/errors/tr/ERR_READ_ERROR
usr/lib/squid/errors/tr/ERR_READ_TIMEOUT
usr/lib/squid/errors/tr/ERR_SECURE_CONNECT_FAIL
usr/lib/squid/errors/tr/ERR_SHUTTING_DOWN
usr/lib/squid/errors/tr/ERR_SOCKET_FAILURE
usr/lib/squid/errors/tr/ERR_TOO_BIG
usr/lib/squid/errors/tr/ERR_UNSUP_HTTPVERSION
usr/lib/squid/errors/tr/ERR_UNSUP_REQ
usr/lib/squid/errors/tr/ERR_URN_RESOLVE
usr/lib/squid/errors/tr/ERR_WRITE_ERROR
usr/lib/squid/errors/tr/ERR_ZERO_SIZE_OBJECT
usr/lib/squid/errors/tr/error-details.txt
usr/local/bin/setddns.pl
var/ipfire/general-functions.pl
var/ipfire/outgoing/bin/outgoingfw.pl
