# config/rootfiles/core/97/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/firewall
etc/modprobe.d/nf_conntrack.conf
srv/web/ipfire/cgi-bin/logs.cgi/firewalllog.dat
srv/web/ipfire/cgi-bin/logs.cgi/firewalllogcountry.dat
srv/web/ipfire/cgi-bin/logs.cgi/firewalllogip.dat
srv/web/ipfire/cgi-bin/logs.cgi/firewalllogport.dat
srv/web/ipfire/cgi-bin/logs.cgi/showrequestfromcountry.dat
srv/web/ipfire/cgi-bin/logs.cgi/showrequestfromip.dat
srv/web/ipfire/cgi-bin/logs.cgi/showrequestfromport.dat
srv/web/ipfire/cgi-bin/optionsfw.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/qos.cgi
srv/web/ipfire/cgi-bin/webaccess.cgi
usr/bin/pgrep
usr/local/bin/qosctrl
usr/local/bin/timectrl
