# config/rootfiles/core/113/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/ntp
etc/rc.d/init.d/unbound
srv/web/ipfire/cgi-bin/logs.cgi/log.dat
var/ipfire/langs
