# config/rootfiles/core/200/filelists/files
etc/rc.d/init.d/unbound
srv/web/ipfire/cgi-bin/ovpnmain.cgi
