# config/rootfiles/core/200/filelists/files
etc/rc.d/init.d/networking/red
etc/rc.d/init.d/unbound
srv/web/ipfire/cgi-bin/dns.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
var/ipfire/backup/exclude
var/ipfire/urlfilter/autoupdate/autoupdate.urls
