# config/rootfiles/oldcore/144/filelists/files
etc/os-release
etc/system-release
etc/issue
srv/web/ipfire/cgi-bin/credits.cgi
var/ipfire/langs
etc/rc.d/init.d/mountkernfs
etc/rc.d/init.d/networking/functions.network
etc/rc.d/init.d/networking/red
