# config/rootfiles/oldcore/95/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/dnsmasq
etc/rc.d/init.d/firewall
etc/rc.d/init.d/networking/red.up/99-geoip-database
lib/udev/network-hotplug-vlan
lib/udev/rules.d/60-net.rules
srv/web/ipfire/cgi-bin/connections.cgi
srv/web/ipfire/cgi-bin/credits.cgi
srv/web/ipfire/cgi-bin/dhcp.cgi
srv/web/ipfire/cgi-bin/firewall.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/logs.cgi/firewalllogcountry.dat
srv/web/ipfire/cgi-bin/mail.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/pppsetup.cgi
srv/web/ipfire/cgi-bin/routing.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
usr/lib/firewall/firewall-lib.pl
usr/lib/firewall/ipsec-block
usr/local/bin/ipsecctrl
usr/local/bin/settime
usr/local/bin/timecheck
var/ipfire/backup/exclude
var/ipfire/dhcpc/dhcpcd.conf
var/ipfire/langs
var/ipfire/network-functions.pl
