# config/rootfiles/oldcore/53/filelists/files
etc/rc.d/init.d/network
etc/rc.d/init.d/static-routes
etc/sysctl.conf
etc/system-release
etc/issue
etc/iproute2/rt_tables
etc/rc.d/init.d/networking/red.up/50-ovpn
etc/rc.d/init.d/networking/red.down/10-ovpn
srv/web/ipfire/cgi-bin/extrahd.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/outgoinggrp.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/pakfire.cgi
srv/web/ipfire/cgi-bin/routing.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
srv/web/ipfire/cgi-bin/logs.cgi/log.dat
var/ipfire/langs/
usr/local/bin/ipsecctrl
usr/local/bin/openvpnctrl
usr/local/bin/vpn-watch
usr/local/bin/rebuildroutes
usr/local/sbin/setup
var/ipfire/main/routing
var/ipfire/menu.d/30-network.menu
opt/pakfire/etc/pakfire.conf
usr/local/bin/backupiso
