# config/rootfiles/core/186/filelists/files
etc/rc.d/init.d/functions
etc/rc.d/init.d/grub-btrfsd
etc/rc.d/rc0.d/K01grub-btrfsd
etc/rc.d/rc3.d/S99grub-btrfsd
etc/rc.d/rc6.d/K01grub-btrfsd
srv/web/ipfire/cgi-bin/vulnerabilities.cgi
