# config/rootfiles/core/156/filelists/files
etc/rc.d/init.d/suricata
etc/sysctl.conf
var/ipfire/graphs.pl
var/ipfire/header.pl
var/ipfire/network-functions.pl
