# config/rootfiles/core/156/filelists/files
etc/rc.d/init.d/suricata
srv/web/ipfire/cgi-bin/firewall.cgi
srv/web/ipfire/cgi-bin/logs.cgi/showrequestfromcountry.dat
srv/web/ipfire/cgi-bin/zoneconf.cgi
var/ipfire/header.pl
var/ipfire/network-functions.pl
