# config/rootfiles/core/156/filelists/files
etc/rc.d/init.d/suricata
etc/sysctl.conf
lib/udev/network-hotplug-bridges
lib/udev/network-hotplug-rename
var/ipfire/graphs.pl
var/ipfire/header.pl
var/ipfire/network-functions.pl
