# config/rootfiles/core/26/files
usr/sbin/pptp
etc/rc.d/init.d/networking/red
srv/web/ipfire/cgi-bin/pppsetup.cgi
