# config/rootfiles/core/80/filelists/files
etc/system-release
etc/issue
etc/logrotate.conf
etc/rc.d/init.d/dhcrelay
etc/rc.d/init.d/dnsmasq
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/routing.cgi
usr/sbin/dhcrelay
var/ipfire/general-functions.pl
var/ipfire/langs
