# config/rootfiles/core/80/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/dhcrelay
usr/local/bin/setddns.pl
usr/sbin/dhcrelay
var/ipfire/general-functions.pl
var/ipfire/langs
