# config/rootfiles/core/100/filelists/files
etc/system-release
etc/issue
etc/dracut.conf
etc/ppp/ip-up
etc/rc.d/init.d/firewall
etc/rc.d/init.d/functions
etc/rc.d/init.d/networking/dhcpcd.exe
etc/rc.d/init.d/ntp
etc/modprobe.d/nf_conntrack.conf
lib/udev/rules.d/60-net.rules
opt/pakfire/etc/pakfire.conf
srv/web/ipfire/cgi-bin/logs.cgi/firewalllog.dat
srv/web/ipfire/cgi-bin/logs.cgi/firewalllogcountry.dat
srv/web/ipfire/cgi-bin/logs.cgi/firewalllogip.dat
srv/web/ipfire/cgi-bin/logs.cgi/firewalllogport.dat
srv/web/ipfire/cgi-bin/logs.cgi/showrequestfromcountry.dat
srv/web/ipfire/cgi-bin/logs.cgi/showrequestfromip.dat
srv/web/ipfire/cgi-bin/logs.cgi/showrequestfromport.dat
srv/web/ipfire/cgi-bin/optionsfw.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/qos.cgi
srv/web/ipfire/cgi-bin/webaccess.cgi
usr/bin/pgrep
usr/local/bin/backupiso
usr/local/bin/qosctrl
usr/local/bin/timectrl
usr/local/bin/wirelessctrl
usr/sbin/firewall-policy
usr/sbin/setup
var/ipfire/dhcpc/dhcpcd-hooks/10-mtu
var/ipfire/langs
