# config/rootfiles/core/189/filelists/files
etc/rc.d/init.d/functions
var/ipfire/header.pl
var/ipfire/ids-functions.pl
