# config/rootfiles/core/69/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/console
var/ipfire/langs/list
var/ipfire/langs/tr.pl
