# config/rootfiles/oldcore/114/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/unbound
srv/web/ipfire/cgi-bin/extrahd.cgi
var/ipfire/langs
