# config/rootfiles/core/38/filelists/files
etc/system-release
etc/issue
etc/e2fsck.conf
etc/httpd/conf/vhosts.d/ipfire-interface.conf
etc/rc.d/init.d/leds
etc/rc.d/init.d/rc
etc/rc.d/init.d/snort
etc/rc.d/init.d/networking/red.up/50-ovpn
etc/rc.d/init.d/networking/dhcpcd.exe
etc/rc.d/rc0.d/K78snort
etc/rc.d/rc0.d/K79leds
etc/rc.d/rc3.d/S21leds
etc/rc.d/rc6.d/K78snort
etc/rc.d/rc6.d/K79leds
etc/udev/rules.d/52-nut-usbups.rules
etc/udev/rules.d/xpp.rules
var/ipfire/backup/bin/backup.pl
var/ipfire/header.pl
var/ipfire/menu.d/50-firewall.menu
var/ipfire/outgoing/bin/outgoingfw.pl
var/ipfire/graphs.pl
var/ipfire/general-functions.pl
etc/modprobe.d/cfg80211
etc/modprobe.d/pcspeaker
var/ipfire/qos/bin/makeqosscripts.pl
srv/web/ipfire/cgi-bin/connections.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/iptables.cgi
srv/web/ipfire/cgi-bin/logs.cgi/config.dat
srv/web/ipfire/cgi-bin/netother.cgi
srv/web/ipfire/cgi-bin/outgoingfw.cgi
srv/web/ipfire/cgi-bin/outgoinggrp.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/services.cgi
srv/web/ipfire/cgi-bin/time.cgi
srv/web/ipfire/cgi-bin/updatexlrator.cgi
srv/web/ipfire/cgi-bin/urlfilter.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
var/ipfire/langs/de.pl
var/ipfire/langs/en.pl
var/ipfire/langs/es.pl
var/ipfire/langs/fr.pl
var/ipfire/langs/list
var/ipfire/main/disable_nf_sip
var/ipfire/outgoing/groups
etc/rc.d/init.d/firewall
etc/rc.d/init.d/ipsec
etc/rc.d/init.d/network
etc/rc.d/init.d/ntp
etc/rc.d/init.d/modules
usr/local/bin/ipsecctrl
usr/local/bin/syslogdctrl
usr/local/bin/wirelessctrl
usr/local/sbin/setup
usr/sbin/openvpn
opt/pakfire/lib/functions.sh
usr/local/share/GeoIP/GeoIP.dat
usr/share/hwdata/pci.ids
usr/share/hwdata/usb.ids
usr/bin/squidGuard
usr/sbin/redirect_wrapper
var/lib/squidguard
var/urlrepo
usr/sbin/klogd
usr/sbin/syslogd
etc/squid/cachemgr.conf
etc/squid/errors
etc/squid/mime.conf
etc/squid/msntauth.conf
etc/squid/squid.conf
srv/web/ipfire/cgi-bin/cachemgr.cgi
usr/bin/squidclient
usr/lib/squid/auth
usr/lib/squid/diskd-daemon
usr/lib/squid/errors.ipfire/English/ERR_ACCESS_DENIED
usr/lib/squid/errors.ipfire/English/ERR_CACHE_ACCESS_DENIED
usr/lib/squid/errors.ipfire/English/ERR_CACHE_MGR_ACCESS_DENIED
usr/lib/squid/errors.ipfire/English/ERR_CANNOT_FORWARD
usr/lib/squid/errors.ipfire/English/ERR_CONNECT_FAIL
usr/lib/squid/errors.ipfire/English/ERR_DNS_FAIL
usr/lib/squid/errors.ipfire/English/ERR_FORWARDING_DENIED
usr/lib/squid/errors.ipfire/English/ERR_FTP_DISABLED
usr/lib/squid/errors.ipfire/English/ERR_FTP_FAILURE
usr/lib/squid/errors.ipfire/English/ERR_FTP_FORBIDDEN
usr/lib/squid/errors.ipfire/English/ERR_FTP_NOT_FOUND
usr/lib/squid/errors.ipfire/English/ERR_FTP_PUT_CREATED
usr/lib/squid/errors.ipfire/English/ERR_FTP_PUT_ERROR
usr/lib/squid/errors.ipfire/English/ERR_FTP_PUT_MODIFIED
usr/lib/squid/errors.ipfire/English/ERR_FTP_UNAVAILABLE
usr/lib/squid/errors.ipfire/English/ERR_INVALID_REQ
usr/lib/squid/errors.ipfire/English/ERR_INVALID_RESP
usr/lib/squid/errors.ipfire/English/ERR_INVALID_URL
usr/lib/squid/errors.ipfire/English/ERR_LIFETIME_EXP
usr/lib/squid/errors.ipfire/English/ERR_NO_RELAY
usr/lib/squid/errors.ipfire/English/ERR_ONLY_IF_CACHED_MISS
usr/lib/squid/errors.ipfire/English/ERR_READ_ERROR
usr/lib/squid/errors.ipfire/English/ERR_READ_TIMEOUT
usr/lib/squid/errors.ipfire/English/ERR_SHUTTING_DOWN
usr/lib/squid/errors.ipfire/English/ERR_SOCKET_FAILURE
usr/lib/squid/errors.ipfire/English/ERR_TOO_BIG
usr/lib/squid/errors.ipfire/English/ERR_UNSUP_REQ
usr/lib/squid/errors.ipfire/English/ERR_URN_RESOLVE
usr/lib/squid/errors.ipfire/English/ERR_WRITE_ERROR
usr/lib/squid/errors.ipfire/English/ERR_ZERO_SIZE_OBJECT
usr/lib/squid/errors.ipfire/German/ERR_ACCESS_DENIED
usr/lib/squid/errors.ipfire/German/ERR_CACHE_ACCESS_DENIED
usr/lib/squid/errors.ipfire/German/ERR_CACHE_MGR_ACCESS_DENIED
usr/lib/squid/errors.ipfire/German/ERR_CANNOT_FORWARD
usr/lib/squid/errors.ipfire/German/ERR_CONNECT_FAIL
usr/lib/squid/errors.ipfire/German/ERR_DNS_FAIL
usr/lib/squid/errors.ipfire/German/ERR_FORWARDING_DENIED
usr/lib/squid/errors.ipfire/German/ERR_FTP_DISABLED
usr/lib/squid/errors.ipfire/German/ERR_FTP_FAILURE
usr/lib/squid/errors.ipfire/German/ERR_FTP_FORBIDDEN
usr/lib/squid/errors.ipfire/German/ERR_FTP_NOT_FOUND
usr/lib/squid/errors.ipfire/German/ERR_FTP_PUT_CREATED
usr/lib/squid/errors.ipfire/German/ERR_FTP_PUT_ERROR
usr/lib/squid/errors.ipfire/German/ERR_FTP_PUT_MODIFIED
usr/lib/squid/errors.ipfire/German/ERR_FTP_UNAVAILABLE
usr/lib/squid/errors.ipfire/German/ERR_INVALID_REQ
usr/lib/squid/errors.ipfire/German/ERR_INVALID_RESP
usr/lib/squid/errors.ipfire/German/ERR_INVALID_URL
usr/lib/squid/errors.ipfire/German/ERR_LIFETIME_EXP
usr/lib/squid/errors.ipfire/German/ERR_NO_RELAY
usr/lib/squid/errors.ipfire/German/ERR_ONLY_IF_CACHED_MISS
usr/lib/squid/errors.ipfire/German/ERR_READ_ERROR
usr/lib/squid/errors.ipfire/German/ERR_READ_TIMEOUT
usr/lib/squid/errors.ipfire/German/ERR_SHUTTING_DOWN
usr/lib/squid/errors.ipfire/German/ERR_SOCKET_FAILURE
usr/lib/squid/errors.ipfire/German/ERR_TOO_BIG
usr/lib/squid/errors.ipfire/German/ERR_UNSUP_REQ
usr/lib/squid/errors.ipfire/German/ERR_URN_RESOLVE
usr/lib/squid/errors.ipfire/German/ERR_WRITE_ERROR
usr/lib/squid/errors.ipfire/German/ERR_ZERO_SIZE_OBJECT
usr/lib/squid/errors/English
usr/lib/squid/errors/German
usr/lib/squid/icons
usr/lib/squid/logfile-daemon
usr/lib/squid/mib.txt
usr/lib/squid/msnt_auth
usr/lib/squid/ncsa_auth
usr/lib/squid/ntlm_auth
usr/lib/squid/smb_auth
usr/lib/squid/smb_auth.pl
usr/lib/squid/smb_auth.sh
usr/lib/squid/squid_ldap_auth
usr/lib/squid/unlinkd
usr/sbin/squid
usr/sbin/updxlrator
var/ipfire/updatexlrator/autocheck/cron.daily
var/ipfire/updatexlrator/autocheck/cron.monthly
var/ipfire/updatexlrator/autocheck/cron.weekly
var/ipfire/updatexlrator/bin/checkdeaddl
var/ipfire/updatexlrator/bin/checkup
var/ipfire/updatexlrator/bin/convert
var/ipfire/updatexlrator/bin/download
var/ipfire/updatexlrator/bin/lscache
var/ipfire/updatexlrator/updxlrator-lib.pl
var/log/cache
var/log/updatexlrator
usr/local/bin/vpn-watch
srv/web/ipfire/html/images/go-down.png
