# config/rootfiles/core/38/filelists/files
etc/system-release
var/ipfire/backup/bin/backup.pl
var/ipfire/header.pl
var/ipfire/menu.d/50-firewall.menu
var/ipfire/outgoing/bin/outgoingfw.pl
var/ipfire/graphs.pl
var/ipfire/general-functions.pl
etc/modprobe.d/cfg80211
etc/modprobe.d/pcspeaker
var/ipfire/qos/bin/makeqosscripts.pl
srv/web/ipfire/cgi-bin/connections.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/iptables.cgi
srv/web/ipfire/cgi-bin/logs.cgi/config.dat
srv/web/ipfire/cgi-bin/netother.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
var/ipfire/langs/de.pl
var/ipfire/langs/en.pl
var/ipfire/langs/fr.pl
etc/rc.d/init.d/firewall
etc/rc.d/init.d/ipsec
etc/rc.d/init.d/network
etc/rc.d/init.d/ntp
etc/rc.d/init.d/modules
usr/local/bin/ipsecctrl
usr/local/bin/syslogdctrl
usr/local/bin/wirelessctrl
usr/sbin/openvpn
opt/pakfire/lib/functions.sh
