# config/rootfiles/core/38/filelists/files
etc/system-release
var/ipfire/backup/bin/backup.pl
var/ipfire/header.pl
etc/modprobe.d/cfg80211
var/ipfire/qos/bin/makeqosscripts.pl
srv/web/ipfire/cgi-bin/connections.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/iptables.cgi
srv/web/ipfire/cgi-bin/logs.cgi/config.dat
srv/web/ipfire/cgi-bin/netother.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
var/ipfire/langs/de.pl
var/ipfire/langs/en.pl
var/ipfire/langs/fr.pl
etc/rc.d/init.d/firewall
etc/rc.d/init.d/ipsec
usr/local/bin/ipsecctrl
usr/local/bin/syslogdctrl
usr/local/bin/wirelessctrl
opt/pakfire/lib/functions.sh

