# config/rootfiles/core/119/filelists/files
etc/system-release
etc/issue
etc/httpd/conf/global.conf
etc/rc.d/init.d/cleanfs
etc/rc.d/init.d/unbound
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
usr/lib/libgcc_s.so.1
var/ipfire/backup/include
var/ipfire/langs
