# config/rootfiles/core/41/filelists/files
etc/system-release
etc/rc.d/init.d/firewall
etc/rc.d/init.d/ntp
srv/web/ipfire/cgi-bin/outgoinggrp.cgi
srv/web/ipfire/cgi-bin/outgoingfw.cgi
srv/web/ipfire/cgi-bin/urlfilter.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
srv/web/ipfire/cgi-bin/speed.cgi
var/ipfire/langs/en.pl
var/ipfire/langs/de.pl
var/ipfire/langs/es.pl
var/ipfire/langs/fr.pl
var/ipfire/outgoing/bin/outgoingfw.pl
usr/local/sbin/setup
lib/openvpn-auth-pam.so
