# config/rootfiles/core/86/filelists/files
etc/system-release
etc/issue
etc/modprobe.d/btmrvl_sdio.conf
etc/rc.d/init.d/leds
opt/pakfire/etc/pakfire.conf
usr/lib/firewall/rules.pl
usr/local/bin/update-bootloader
var/ipfire/header.pl
var/ipfire/langs
var/ipfire/lang.pl
