# config/rootfiles/core/17/files
etc/squid
etc/init.d/connectd
usr/lib/squid
usr/sbin/squid
