# config/rootfiles/oldcore/179/filelists/files
etc/rc.d/init.d/networking/red
lib/udev/rules.d/61-extrahd.rules
srv/web/ipfire/cgi-bin/extrahd.cgi
srv/web/ipfire/cgi-bin/fwhosts.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/vulnerabilities.cgi
var/ipfire/extrahd/bin/extrahd.pl
