# config/rootfiles/core/169/filelists/files
etc/httpd/conf/vhosts.d/ipfire-interface.conf
etc/httpd/conf/vhosts.d/ipfire-interface-ssl.conf
etc/ntp.conf
etc/rc.d/helper/aws-setup
etc/rc.d/helper/azure-setup
etc/rc.d/helper/exoscale-setup
etc/rc.d/helper/gcp-setup
etc/rc.d/helper/oci-setup
etc/rc.d/init.d/localnet
etc/rc.d/init.d/mountkernfs
etc/rc.d/init.d/ntp
etc/rc.d/init.d/udev
etc/sysctl.conf
lib/firmware/amdgpu/aldebaran_mec2.bin
lib/firmware/amdgpu/aldebaran_mec.bin
lib/firmware/amdgpu/aldebaran_rlc.bin
lib/firmware/amdgpu/aldebaran_smc.bin
lib/firmware/amdgpu/aldebaran_sos.bin
lib/firmware/amdgpu/aldebaran_ta.bin
lib/firmware/amdgpu/beige_goby_ce.bin
lib/firmware/amdgpu/beige_goby_dmcub.bin
lib/firmware/amdgpu/beige_goby_me.bin
lib/firmware/amdgpu/beige_goby_mec2.bin
lib/firmware/amdgpu/beige_goby_mec.bin
lib/firmware/amdgpu/beige_goby_pfp.bin
lib/firmware/amdgpu/beige_goby_smc.bin
lib/firmware/amdgpu/beige_goby_ta.bin
lib/firmware/amdgpu/beige_goby_vcn.bin
lib/firmware/amdgpu/dcn_3_1_6_dmcub.bin
lib/firmware/amdgpu/dimgrey_cavefish_ce.bin
lib/firmware/amdgpu/dimgrey_cavefish_dmcub.bin
lib/firmware/amdgpu/dimgrey_cavefish_me.bin
lib/firmware/amdgpu/dimgrey_cavefish_mec2.bin
lib/firmware/amdgpu/dimgrey_cavefish_mec.bin
lib/firmware/amdgpu/dimgrey_cavefish_pfp.bin
lib/firmware/amdgpu/dimgrey_cavefish_smc.bin
lib/firmware/amdgpu/dimgrey_cavefish_sos.bin
lib/firmware/amdgpu/dimgrey_cavefish_ta.bin
lib/firmware/amdgpu/dimgrey_cavefish_vcn.bin
lib/firmware/amdgpu/gc_10_3_7_rlc.bin
lib/firmware/amdgpu/green_sardine_asd.bin
lib/firmware/amdgpu/green_sardine_dmcub.bin
lib/firmware/amdgpu/green_sardine_ta.bin
lib/firmware/amdgpu/navi10_asd.bin
lib/firmware/amdgpu/navi10_ce.bin
lib/firmware/amdgpu/navi10_me.bin
lib/firmware/amdgpu/navi10_mec2.bin
lib/firmware/amdgpu/navi10_mec.bin
lib/firmware/amdgpu/navi10_pfp.bin
lib/firmware/amdgpu/navi10_sos.bin
lib/firmware/amdgpu/navi10_ta.bin
lib/firmware/amdgpu/navi12_asd.bin
lib/firmware/amdgpu/navi12_ce.bin
lib/firmware/amdgpu/navi12_me.bin
lib/firmware/amdgpu/navi12_mec2.bin
lib/firmware/amdgpu/navi12_mec.bin
lib/firmware/amdgpu/navi12_pfp.bin
lib/firmware/amdgpu/navi12_rlc.bin
lib/firmware/amdgpu/navi12_sos.bin
lib/firmware/amdgpu/navi12_ta.bin
lib/firmware/amdgpu/navi14_asd.bin
lib/firmware/amdgpu/navi14_ce.bin
lib/firmware/amdgpu/navi14_me.bin
lib/firmware/amdgpu/navi14_mec2.bin
lib/firmware/amdgpu/navi14_mec.bin
lib/firmware/amdgpu/navi14_pfp.bin
lib/firmware/amdgpu/navi14_sos.bin
lib/firmware/amdgpu/navi14_ta.bin
lib/firmware/amdgpu/navy_flounder_ce.bin
lib/firmware/amdgpu/navy_flounder_dmcub.bin
lib/firmware/amdgpu/navy_flounder_me.bin
lib/firmware/amdgpu/navy_flounder_mec2.bin
lib/firmware/amdgpu/navy_flounder_mec.bin
lib/firmware/amdgpu/navy_flounder_pfp.bin
lib/firmware/amdgpu/navy_flounder_rlc.bin
lib/firmware/amdgpu/navy_flounder_sdma.bin
lib/firmware/amdgpu/navy_flounder_ta.bin
lib/firmware/amdgpu/navy_flounder_vcn.bin
lib/firmware/amdgpu/picasso_asd.bin
lib/firmware/amdgpu/picasso_ta.bin
lib/firmware/amdgpu/psp_13_0_8_asd.bin
lib/firmware/amdgpu/psp_13_0_8_ta.bin
lib/firmware/amdgpu/raven2_asd.bin
lib/firmware/amdgpu/raven2_ta.bin
lib/firmware/amdgpu/raven_asd.bin
lib/firmware/amdgpu/raven_ta.bin
lib/firmware/amdgpu/renoir_asd.bin
lib/firmware/amdgpu/renoir_dmcub.bin
lib/firmware/amdgpu/renoir_ta.bin
lib/firmware/amdgpu/sienna_cichlid_ce.bin
lib/firmware/amdgpu/sienna_cichlid_dmcub.bin
lib/firmware/amdgpu/sienna_cichlid_me.bin
lib/firmware/amdgpu/sienna_cichlid_mec2.bin
lib/firmware/amdgpu/sienna_cichlid_mec.bin
lib/firmware/amdgpu/sienna_cichlid_pfp.bin
lib/firmware/amdgpu/sienna_cichlid_sos.bin
lib/firmware/amdgpu/sienna_cichlid_ta.bin
lib/firmware/amdgpu/sienna_cichlid_vcn.bin
lib/firmware/amdgpu/vangogh_asd.bin
lib/firmware/amdgpu/vangogh_ce.bin
lib/firmware/amdgpu/vangogh_me.bin
lib/firmware/amdgpu/vangogh_mec2.bin
lib/firmware/amdgpu/vangogh_mec.bin
lib/firmware/amdgpu/vangogh_pfp.bin
lib/firmware/amdgpu/vangogh_rlc.bin
lib/firmware/amdgpu/vega10_asd.bin
lib/firmware/amdgpu/vega12_asd.bin
lib/firmware/amdgpu/vega20_asd.bin
lib/firmware/amdgpu/yellow_carp_asd.bin
lib/firmware/amdgpu/yellow_carp_ce.bin
lib/firmware/amdgpu/yellow_carp_dmcub.bin
lib/firmware/amdgpu/yellow_carp_me.bin
lib/firmware/amdgpu/yellow_carp_mec2.bin
lib/firmware/amdgpu/yellow_carp_mec.bin
lib/firmware/amdgpu/yellow_carp_pfp.bin
lib/firmware/amdgpu/yellow_carp_ta.bin
lib/firmware/ath10k/QCA4019/hw1.0/board-2.bin
lib/firmware/ath10k/QCA9888/hw2.0/board-2.bin
lib/firmware/ath10k/QCA9888/hw2.0/firmware-5.bin
lib/firmware/ath10k/QCA9984/hw1.0/firmware-5.bin
lib/firmware/ath10k/QCA99X0/hw2.0/board-2.bin
lib/firmware/ath10k/WCN3990/hw1.0/board-2.bin
lib/firmware/ath11k/QCN9074/hw1.0/amss.bin
lib/firmware/ath11k/QCN9074/hw1.0/board-2.bin
lib/firmware/ath11k/QCN9074/hw1.0/m3.bin
lib/firmware/ath11k/WCN6750/hw1.0/board-2.bin
lib/firmware/ath11k/WCN6750/hw1.0/wpss.b00
lib/firmware/ath11k/WCN6750/hw1.0/wpss.b01
lib/firmware/ath11k/WCN6750/hw1.0/wpss.b02
lib/firmware/ath11k/WCN6750/hw1.0/wpss.b03
lib/firmware/ath11k/WCN6750/hw1.0/wpss.b04
lib/firmware/ath11k/WCN6750/hw1.0/wpss.b05
lib/firmware/ath11k/WCN6750/hw1.0/wpss.b06
lib/firmware/ath11k/WCN6750/hw1.0/wpss.b07
lib/firmware/ath11k/WCN6750/hw1.0/wpss.b08
lib/firmware/ath11k/WCN6750/hw1.0/wpss.mdt
lib/firmware/ath11k/WCN6855/hw2.0/amss.bin
lib/firmware/atmel/wilc1000_wifi_firmware-1.bin
lib/firmware/brcm/brcmfmac43430-sdio.beagle,beaglev-starlight-jh7100-a1.txt
lib/firmware/brcm/brcmfmac43430-sdio.beagle,beaglev-starlight-jh7100-r0.txt
lib/firmware/brcm/brcmfmac43430-sdio.starfive,visionfive-v1.txt
lib/firmware/brcm/brcmfmac43455-sdio.pine64,quartz64-b.txt
lib/firmware/i915/adlp_guc_70.1.1.bin
lib/firmware/i915/bxt_guc_70.1.1.bin
lib/firmware/i915/cml_guc_70.1.1.bin
lib/firmware/i915/dg1_guc_70.1.1.bin
lib/firmware/i915/dg2_dmc_ver2_06.bin
lib/firmware/i915/dg2_guc_70.1.2.bin
lib/firmware/i915/ehl_guc_70.1.1.bin
lib/firmware/i915/glk_guc_70.1.1.bin
lib/firmware/i915/icl_guc_70.1.1.bin
lib/firmware/i915/kbl_guc_70.1.1.bin
lib/firmware/i915/skl_guc_70.1.1.bin
lib/firmware/i915/tgl_guc_70.1.1.bin
lib/firmware/intel/ibt-0040-0041.sfi
lib/firmware/intel/ibt-0040-1020.sfi
lib/firmware/intel/ibt-0040-2120.sfi
lib/firmware/intel/ibt-0040-4150.sfi
lib/firmware/intel/ibt-0041-0041.sfi
lib/firmware/intel/ibt-1040-0041.sfi
lib/firmware/intel/ibt-1040-1020.sfi
lib/firmware/intel/ibt-1040-2120.sfi
lib/firmware/intel/ibt-1040-4150.sfi
lib/firmware/intel/ibt-19-0-0.sfi
lib/firmware/intel/ibt-19-0-1.sfi
lib/firmware/intel/ibt-19-0-4.sfi
lib/firmware/intel/ibt-19-16-4.sfi
lib/firmware/intel/ibt-19-240-1.sfi
lib/firmware/intel/ibt-19-240-4.sfi
lib/firmware/intel/ibt-19-32-0.sfi
lib/firmware/intel/ibt-19-32-1.sfi
lib/firmware/intel/ibt-19-32-4.sfi
lib/firmware/intel/ibt-20-0-3.sfi
lib/firmware/intel/ibt-20-1-3.sfi
lib/firmware/intel/ibt-20-1-4.sfi
lib/firmware/intel/ice/ddp/ice-1.3.28.0.pkg
lib/firmware/iwlwifi-9000-pu-b0-jf-b0-46.ucode
lib/firmware/iwlwifi-9260-th-b0-jf-b0-46.ucode
lib/firmware/iwlwifi-cc-a0-73.ucode
lib/firmware/iwlwifi-Qu-b0-hr-b0-73.ucode
lib/firmware/iwlwifi-Qu-b0-jf-b0-73.ucode
lib/firmware/iwlwifi-Qu-c0-hr-b0-73.ucode
lib/firmware/iwlwifi-Qu-c0-jf-b0-73.ucode
lib/firmware/iwlwifi-QuZ-a0-hr-b0-73.ucode
lib/firmware/iwlwifi-QuZ-a0-jf-b0-73.ucode
lib/firmware/iwlwifi-so-a0-gf4-a0-73.ucode
lib/firmware/iwlwifi-so-a0-gf4-a0.pnvm
lib/firmware/iwlwifi-so-a0-gf-a0-73.ucode
lib/firmware/iwlwifi-so-a0-gf-a0.pnvm
lib/firmware/iwlwifi-so-a0-hr-b0-73.ucode
lib/firmware/iwlwifi-so-a0-jf-b0-73.ucode
lib/firmware/iwlwifi-ty-a0-gf-a0-73.ucode
lib/firmware/iwlwifi-ty-a0-gf-a0.pnvm
lib/firmware/mediatek/BT_RAM_CODE_MT7922_1_1_hdr.bin
lib/firmware/mediatek/mt8183/scp.img
lib/firmware/mediatek/mt8192/scp.img
lib/firmware/mediatek/mt8195/scp.img
lib/firmware/mediatek/WIFI_MT7922_patch_mcu_1_1_hdr.bin
lib/firmware/mediatek/WIFI_RAM_CODE_MT7922_1.bin
lib/firmware/qcom/a330_pfp.fw
lib/firmware/qcom/a330_pm4.fw
lib/firmware/qcom/a420_pfp.fw
lib/firmware/qcom/a420_pm4.fw
lib/firmware/qcom/apq8096/adsp.mbn
lib/firmware/qcom/apq8096/adspr.jsn
lib/firmware/qcom/apq8096/adspua.jsn
lib/firmware/qcom/apq8096/mba.mbn
lib/firmware/qcom/apq8096/modem.mbn
lib/firmware/qcom/apq8096/modemr.jsn
lib/firmware/qcom/leia_pfp_470.fw
lib/firmware/qcom/leia_pm4_470.fw
lib/firmware/rtl_bt/rtl8852au_fw.bin
lib/firmware/rtl_bt/rtl8852cu_config.bin
lib/firmware/rtl_bt/rtl8852cu_fw.bin
lib/firmware/rtw89/rtw8852c_fw.bin
opt/pakfire/etc/pakfire.conf
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/html/images/qr-code.png
srv/web/ipfire/html/images/qr-code.svg
srv/web/ipfire/html/themes/ipfire/include/functions.pl
usr/sbin/setup
var/ipfire/header.pl
var/ipfire/main/manualpages
