# config/rootfiles/core/87/filelists/files
etc/system-release
etc/issue
etc/modprobe.d/btmrvl_sdio.conf
etc/rc.d/init.d/collectd
etc/rc.d/init.d/leds
opt/pakfire/etc/pakfire.conf
usr/lib/firewall/rules.pl
usr/local/bin/backupiso
usr/local/bin/update-bootloader
usr/local/share/GeoIP/GeoIP.dat
var/ipfire/header.pl
var/ipfire/langs
var/ipfire/lang.pl
var/ipfire/ovpn/openssl/ovpn.cnf
