# config/rootfiles/core/7/files
etc/collectd.conf
etc/init.d/collectd
srv/web/ipfire/cgi-bin/graphs.cgi
srv/web/ipfire/cgi-bin/hardwaregraphs.cgi
srv/web/ipfire/cgi-bin/network.cgi
usr/lib/collectd
usr/local/bin/makegraphs
usr/sbin/collectd
var/ipfire/graphs.pl
var/ipfire/menu.d/20-status.menu


