# config/rootfiles/core/test/filelists/files
etc/system-release
etc/rc.d/init.d/checkfstab
usr/local/bin/rebuild-initrd
srv/web/ipfire/cgi-bin/speed.cgi
