# config/rootfiles/core/56/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/smartenabler
srv/web/ipfire/cgi-bin/hardwaregraphs.cgi
srv/web/ipfire/cgi-bin/media.cgi
usr/local/bin/hddshutdown
usr/local/bin/makegraphs
