# config/rootfiles/core/56/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/mountfs
etc/rc.d/init.d/smartenabler
etc/rc.d/init.d/networking/red.up/10-static-routes
etc/rc.d/init.d/networking/red.down/10-static-routes
srv/web/ipfire/cgi-bin/chpasswd.cgi
srv/web/ipfire/cgi-bin/hardwaregraphs.cgi
srv/web/ipfire/cgi-bin/media.cgi
usr/local/bin/hddshutdown
usr/local/bin/makegraphs
usr/local/bin/openvpnctrl
usr/local/bin/readhash
