# config/rootfiles/core/123/filelists/files
etc/system-release
etc/issue
etc/rc.d/helper/aws-setup
etc/rc.d/init.d/aws
etc/rc.d/init.d/collectd
etc/rc.d/init.d/unbound
etc/rc.d/rcsysinit.d/S74aws
etc/sysctl.conf
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/logs.cgi/ids.dat
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
usr/local/bin/backupiso
usr/lib/libssp.so.0
usr/lib/libssp.so.0.0.0
usr/sbin/dhclient
usr/share/GeoIP/GeoIP.dat
var/ipfire/backup/bin/backup.pl
var/ipfire/backup/exclude
var/ipfire/langs
var/ipfire/aws-functions.pl
var/ipfire/header.pl
var/ipfire/general-functions.pl
