# config/rootfiles/core/136/filelists/files
etc/system-release
etc/issue
srv/web/ipfire/cgi-bin/credits.cgi
var/ipfire/langs
etc/rc.d/init.d/firewall
usr/local/bin/xt_geoip_update
srv/web/ipfire/cgi-bin/logs.cgi/log.dat
