# config/rootfiles/core/161/filelists/files
etc/rc.d/init.d/firewall
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/qos.cgi
usr/bin/hexdump
usr/lib/firewall/rules.pl
usr/libexec/ipsec/_updown
usr/local/bin/hddshutdown
usr/local/bin/makegraphs
var/ipfire/qos/bin/makeqosscripts.pl
