# config/rootfiles/core/161/filelists/files
etc/rc.d/init.d/firewall
etc/rc.d/init.d/partresize
opt/pakfire/etc/pakfire.conf
opt/pakfire/lib/functions.pl
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/qos.cgi
usr/bin/2to3
usr/bin/hexdump
usr/lib/firewall/rules.pl
usr/libexec/ipsec/_updown
usr/local/bin/hddshutdown
usr/local/bin/makegraphs
var/ipfire/general-functions.pl
var/ipfire/qos/bin/makeqosscripts.pl
