# config/rootfiles/core/161/filelists/files
etc/rc.d/init.d/firewall
etc/rc.d/init.d/networking/red
etc/rc.d/init.d/partresize
etc/rc.d/init.d/suricata
etc/rc.d/helper/azure-setup
etc/suricata/suricata.yaml
opt/pakfire/etc/pakfire.conf
opt/pakfire/lib/functions.pl
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/pakfire.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/qos.cgi
srv/web/ipfire/cgi-bin/speed.cgi
usr/bin/2to3
usr/bin/hexdump
usr/lib/firewall/rules.pl
usr/local/bin/hddshutdown
usr/local/bin/makegraphs
var/ipfire/backup/exclude
var/ipfire/general-functions.pl
var/ipfire/qos/bin/makeqosscripts.pl
