# config/rootfiles/oldcore/129/filelists/files
etc/system-release
etc/issue
var/ipfire/langs
etc/rc.d/init.d/firewall
etc/rc.d/init.d/network
etc/rc.d/init.d/networking/red.up/50-ipsec
etc/rc.d/init.d/unbound
srv/web/ipfire/cgi-bin/credits.cgi
srv/web/ipfire/cgi-bin/dhcp.cgi
srv/web/ipfire/cgi-bin/dnsforward.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/netovpnsrv.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
usr/bin/readelf
usr/bin/strings
usr/lib/firewall/firewall-lib.pl
usr/lib/firewall/ipsec-policy
usr/local/bin/ipsec-interfaces
usr/local/bin/ipsecctrl
