# config/rootfiles/core/170/filelists/files
etc/rc.d/init.d/firewall
opt/pakfire/lib/functions.pl
opt/pakfire/pakfire
srv/web/ipfire/cgi-bin/aliases.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/ipblocklist.cgi
srv/web/ipfire/cgi-bin/pakfire.cgi
srv/web/ipfire/cgi-bin/services.cgi
srv/web/ipfire/cgi-bin/logs.cgi/ids.dat
srv/web/ipfire/cgi-bin/logs.cgi/ipblocklists.dat
srv/web/ipfire/cgi-bin/logs.cgi/log.dat
srv/web/ipfire/cgi-bin/logs.cgi/showrequestfromblocklist.dat
srv/web/ipfire/html/themes/ipfire/include/functions.pl
usr/lib/firewall/rules.pl
usr/local/bin/setaliases
usr/local/bin/update-ipblocklists
usr/share/terminfo/t/tmux
usr/share/terminfo/t/tmux-256color
usr/share/terminfo/t/tmux-direct
var/ipfire/backup/bin/backup.pl
var/ipfire/backup/include
var/ipfire/ipblocklist-functions.pl
var/ipfire/menu.d/50-firewall.menu
var/ipfire/menu.d/70-log.menu
var/ipfire/network-functions.pl
var/ipfire/urlfilter/autoupdate/autoupdate.urls
var/spool/cron/root.orig
