# config/rootfiles/core/34/filelists/files
etc/system-release
etc/rc.d/init.d/networking/red
etc/rc.d/init.d/tmpfs
boot/grub/grub.conf
