# config/udev/90-hwrng.rules
# do not edit this file, it will be overwritten on update

SUBSYSTEM!="tty", GOTO="ipfire_hwrng_end"
KERNEL!="ttyACM[0-9]*", GOTO="ipfire_hwrng_end"

SUBSYSTEMS=="usb-serial", ENV{.ID_PORT}="$attr{port_number}"

IMPORT{builtin}="usb_id"
ENV{ID_SERIAL}=="", GOTO="ipfire_hwrng_end"
SUBSYSTEMS=="usb", ENV{ID_USB_INTERFACE_NUM}="$attr{bInterfaceNumber}"
ENV{ID_USB_INTERFACE_NUM}=="", GOTO="ipfire_hwrng_end"

ATTRS{manufacturer}=="IPFire.org", ATTRS{product}=="Random Number Generator*", ENV{.ID_PORT}=="", SYMLINK+="hwrngtty" RUN+="/bin/stty raw -echo -ixoff -F /dev/hwrngtty 115200" RUN+="/etc/rc.d/init.d/rngd udev-event"

LABEL="ipfire_hwrng_end"
ACTION=="add|remove", KERNEL=="hw_random", RUN+="/etc/rc.d/init.d/rngd udev-event"
