# config/udev/99-aqm.rules
# Call the AQM script
SUBSYSTEM=="net", RUN+="/usr/lib/udev/network-aqm"
