# config/udev/99-codel.rules
# Call the enable codel script.
SUBSYSTEM=="net", RUN+="/lib/udev/enable_codel"
