# config/udev/99-offloading.rules
# Call the offloading script
SUBSYSTEM=="net", RUN+="/lib/udev/network-offloading"
