# config/profile.d/glib2-locale.sh
# Use the current locale charset for filenames
# in applications using GLib
export G_FILENAME_ENCODING=@locale
