# config/rootfiles/core/145/filelists/files
etc/os-release
etc/system-release
etc/issue
srv/web/ipfire/cgi-bin/credits.cgi
var/ipfire/langs
etc/rc.d/init.d/firewall
etc/rc.d/init.d/vnstat
opt/pakfire/lib/functions.pl
srv/web/ipfire/cgi-bin/logs.cgi/ovpnclients.dat
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/system.cgi
srv/web/ipfire/cgi-bin/traffic.cgi
usr/local/bin/makegraphs
usr/sbin/openvpn-metrics
var/ipfire/menu.d/70-log.menu
var/ipfire/general-functions.pl
var/ipfire/graphs.pl
