# config/rootfiles/core/145/filelists/files
etc/os-release
etc/system-release
etc/issue
srv/web/ipfire/cgi-bin/credits.cgi
var/ipfire/langs
etc/rc.d/init.d/firewall
etc/rc.d/init.d/vnstat
opt/pakfire/lib/functions.pl
usr/local/bin/makegraphs
srv/web/ipfire/cgi-bin/system.cgi
srv/web/ipfire/cgi-bin/traffic.cgi
