# config/rootfiles/core/145/filelists/files
etc/os-release
etc/system-release
etc/issue
srv/web/ipfire/cgi-bin/credits.cgi
var/ipfire/langs
etc/rc.d/init.d/firewall
etc/rc.d/init.d/vnstat
opt/pakfire/lib/functions.pl
srv/web/ipfire/cgi-bin/logs.cgi/ovpnclients.dat
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/optionsfw.cgi
srv/web/ipfire/cgi-bin/system.cgi
srv/web/ipfire/cgi-bin/traffic.cgi
usr/lib/firewall/rules.pl
usr/local/bin/makegraphs
usr/sbin/openvpn-metrics
usr/share/hwdata/pci.ids
var/ipfire/menu.d/70-log.menu
var/ipfire/ids-functions.pl
var/ipfire/general-functions.pl
var/ipfire/graphs.pl
usr/bin/squidGuard
usr/lib/sasl2/libsasldb.so.3.0.0
usr/lib/perl5/5.30.0/xxxMACHINExxx-linux-thread-multi/auto/DB_File/DB_File.so
usr/lib/perl5/site_perl/5.30.0/xxxMACHINExxx-linux-thread-multi/auto/BerkeleyDB/BerkeleyDB.so
usr/lib/python2.7/lib-dynload/_bsddb.so
usr/sbin/sasldblistusers2
usr/sbin/saslpasswd2
usr/sbin/pluginviewer
usr/bin/rrdcached
usr/bin/pango-view
usr/lib/librrd.so.8.2.1
usr/lib/libpoppler.so.66.0.0
usr/lib/libgd.so.2.0.0
usr/lib/libcairo-gobject.so.2.11600.0
usr/lib/collectd/rrdtool.so
usr/lib/collectd/rrdcached.so
usr/lib/libcairo.so.2.11600.0
usr/lib/perl5/site_perl/5.30.0/xxxMACHINExxx-linux-thread-multi/auto/GD/GD.so
usr/lib/libcairo-script-interpreter.so.2.11600.0
usr/lib/libpangocairo-1.0.so.0.3000.1
