# config/rootfiles/core/138/filelists/files
etc/system-release
etc/issue
srv/web/ipfire/cgi-bin/credits.cgi
var/ipfire/langs
etc/httpd/conf/vhosts.d/ipfire-interface.conf
etc/httpd/conf/vhosts.d/ipfire-interface-ssl.conf
etc/rc.d/init.d/functions
etc/rc.d/init.d/networking/red.up/23-suricata
etc/rc.d/init.d/unbound
etc/suricata/suricata.yaml
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/mail.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
usr/sbin/convert-snort
usr/lib/firewall/firewall-lib.pl
var/ipfire/ids-functions.pl
