# config/rootfiles/core/138/filelists/files
etc/system-release
etc/issue
srv/web/ipfire/cgi-bin/credits.cgi
var/ipfire/langs
etc/httpd/conf/vhosts.d/interface-ipfire.conf
etc/httpd/conf/vhosts.d/interface-ipfire-ssl.conf
etc/rc.d/init.d/functions
etc/rc.d/init.d/unbound
srv/web/ipfire/cgi-bin/mail.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
usr/lib/firewall/firewall-lib.pl
