# config/rootfiles/core/143/filelists/files
etc/system-release
etc/issue
srv/web/ipfire/cgi-bin/credits.cgi
var/ipfire/langs
etc/rc.d/init.d/firewall
etc/rc.d/init.d/localnet
srv/web/ipfire/cgi-bin/dhcp.cgi
srv/web/ipfire/cgi-bin/fireinfo.cgi
srv/web/ipfire/cgi-bin/netother.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
usr/lib/firewall/rules.pl
var/ipfire/backup/bin/backup.pl
var/ipfire/suricata/oinkmaster.conf
