# config/rootfiles/core/143/filelists/files
etc/system-release
etc/issue
srv/web/ipfire/cgi-bin/credits.cgi
var/ipfire/langs
etc/rc.d/init.d/firewall
etc/rc.d/init.d/localnet
etc/rc.d/init.d/unbound
etc/suricata/suricata.yaml
srv/web/ipfire/cgi-bin/dhcp.cgi
srv/web/ipfire/cgi-bin/fireinfo.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/mail.cgi
srv/web/ipfire/cgi-bin/netother.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
srv/web/ipfire/cgi-bin/zoneconf.cgi
usr/lib/firewall/rules.pl
usr/sbin/convert-snort
var/ipfire/backup/bin/backup.pl
var/ipfire/backup/include
var/ipfire/ids-functions.pl
var/ipfire/suricata/oinkmaster.conf
var/ipfire/suricata/ruleset-sources
