# config/rootfiles/core/36/filelists/files
etc/system-release
etc/rc.d/init.d/firewall
