# config/rootfiles/oldcore/29/filelists/files
etc/collectd.conf
etc/rc.d/init.d/networking/red
etc/system-release
lib/firmware/av7110/bootcode.bin
srv/web/ipfire/cgi-bin/aliases.cgi
srv/web/ipfire/cgi-bin/credits.cgi
srv/web/ipfire/cgi-bin/backup.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/html/blob.gif
srv/web/ipfire/html/images/info.gif
srv/web/ipfire/html/images/openvpn.png
usr/bin/nsupdate
var/ipfire/backup/include
var/ipfire/langs/list
var/ipfire/langs/de.pl
var/ipfire/langs/en.pl
var/ipfire/mac/settings
var/ipfire/dns/settings
var/ipfire/menu.d/50-firewall.menu
var/ipfire/menu.d/30-network.menu
srv/web/ipfire/cgi-bin/mac.cgi
srv/web/ipfire/cgi-bin/dns.cgi
srv/web/ipfire/html/images/indicator.gif
srv/web/ipfire/html/images/dns_link.png
srv/web/ipfire/cgi-bin/qos.cgi
srv/web/ipfire/cgi-bin/index.cgi
