# config/rootfiles/core/71/filelists/files
etc/system-release
etc/issue
var/ipfire/general-functions.pl
var/ipfire/header.pl
var/ipfire/menu.d/10-system.menu
etc/rc.d/init.d/wlanclient
etc/rc.d/init.d/networking/functions.network
etc/rc.d/init.d/networking/wpa_supplicant.exe
etc/rc.d/init.d/networking/red
etc/rc.d/rc0.d/K82wlanclient
etc/rc.d/rc3.d/S19wlanclient
etc/rc.d/rc6.d/K82wlanclient
usr/local/bin/wirelessclient
srv/web/ipfire/cgi-bin/wirelessclient.cgi
var/ipfire/langs
var/ipfire/backup/include
