# config/rootfiles/core/178/filelists/files
etc/rc.d/init.d/networking/red
srv/web/ipfire/cgi-bin/extrahd.cgi
srv/web/ipfire/cgi-bin/fwhosts.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
var/ipfire/extrahd/bin/extrahd.pl
