# config/rootfiles/oldcore/139/filelists/files
etc/system-release
etc/issue
srv/web/ipfire/cgi-bin/credits.cgi
var/ipfire/langs
etc/httpd/conf/vhosts.d/ipfire-interface.conf
etc/httpd/conf/vhosts.d/ipfire-interface-ssl.conf
etc/rc.d/init.d/functions
etc/rc.d/init.d/leds
etc/rc.d/init.d/networking/dhcpcd.exe
etc/rc.d/init.d/networking/red
etc/rc.d/init.d/networking/functions.network
etc/rc.d/init.d/networking/red.down/99-beep
etc/rc.d/init.d/networking/red.up/23-suricata
etc/rc.d/init.d/networking/red.up/99-beep
etc/rc.d/init.d/unbound
etc/suricata/suricata.yaml
etc/ppp/ip-down
etc/ppp/ip-up
srv/web/ipfire/cgi-bin/captive.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/mail.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
usr/bin/pcregrep
usr/sbin/convert-snort
usr/lib/firewall/firewall-lib.pl
var/ipfire/dhcpc/dhcpcd-hooks/10-mtu
var/ipfire/ids-functions.pl
var/ipfire/proxy/calamaris/bin/mkreport
