# config/rootfiles/oldcore/164/filelists/files
etc/collectd.conf
etc/rc.d/init.d/firewall
etc/rc.d/init.d/leds
etc/rc.d/init.d/partresize
etc/rc.d/init.d/squid
etc/rc.d/init.d/unbound
etc/suricata/suricata.yaml
etc/sysctl.conf
opt/pakfire/lib/functions.pl
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/logs.cgi/log.dat
srv/web/ipfire/cgi-bin/optionsfw.cgi
srv/web/ipfire/cgi-bin/pakfire.cgi
srv/web/ipfire/html/include/pakfire.js
usr/local/bin/update-ids-ruleset
usr/sbin/convert-ids-multiple-providers
usr/sbin/convert-snort
var/ipfire/backup/bin/backup.pl
var/ipfire/backup/include
var/ipfire/graphs.pl
var/ipfire/ids-functions.pl
var/ipfire/urlfilter/autoupdate/autoupdate.urls
