# config/rootfiles/core/173/filelists/files
etc/ppp/ip-up
etc/rc.d/init.d/networking/dhcpcd.exe
etc/rc.d/init.d/networking/functions.network
etc/rc.d/init.d/networking/red
lib/udev/network-hotplug-rename
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/pppsetup.cgi
var/ipfire/dhcpc/dhcpcd.conf
