# config/rootfiles/core/39/filelists/files
etc/system-release
etc/init.d/networking/red
srv/web/ipfire/cgi-bin/pppsetup.cgi
