# config/rootfiles/core/39/filelists/files
etc/system-release
etc/rc.d/init.d/dnsmasq
etc/rc.d/init.d/networking/red
etc/rc.d/helper/getdnsfromdhcpc.pl
srv/web/ipfire/cgi-bin/pppsetup.cgi
usr/local/bin/backupiso
usr/local/bin/updxlratorctrl
usr/local/bin/urlfilterctrl
srv/web/ipfire/cgi-bin/updatexlrator.cgi
srv/web/ipfire/cgi-bin/urlfilter.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
/opt/pakfire/lib/functions.pl
var/ipfire/langs/de.pl
var/ipfire/langs/en.pl
var/ipfire/langs/es.pl
var/ipfire/langs/fr.pl