# config/rootfiles/core/65/filelists/files
etc/system-release
etc/issue
etc/sysctl.conf
etc/securetty
opt/pakfire/etc/pakfire.conf
etc/rc.d/init.d/firewall
etc/rc.d/init.d/fsresize
etc/rc.d/init.d/halt
etc/rc.d/init.d/mountfs
etc/rc.d/init.d/network
etc/rc.d/init.d/partresize
etc/rc.d/init.d/reboot
etc/rc.d/init.d/sendsignals
etc/rc.d/init.d/setclock
etc/rc.d/init.d/tmpfs
etc/rc.d/init.d/udev
etc/rc.d/init.d/udev_retry
etc/rc.d/rc0.d/K47setclock
etc/rc.d/rc6.d/K47setclock
srv/web/ipfire/cgi-bin/credits.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/qos.cgi
srv/web/ipfire/cgi-bin/services.cgi
srv/web/ipfire/cgi-bin/updatexlrator.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
var/ipfire/langs
var/ipfire/outgoing/bin/outgoingfw.pl
var/ipfire/qos/bin/makeqosscripts.pl
usr/local/bin/readhash
usr/local/bin/scanhd
