# config/rootfiles/core/45/filelists/files
etc/system-release
etc/issue
etc/httpd/conf/ssl-global.conf
etc/rc.d/init.d/functions
srv/web/ipfire/cgi-bin/connections.cgi
srv/web/ipfire/cgi-bin/fireinfo.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/pakfire.cgi
srv/web/ipfire/cgi-bin/pppsetup.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
usr/sbin/updxlrator
var/ipfire/outgoing/bin/outgoingfw.pl
srv/web/ipfire/cgi-bin/logs.cgi/firewalllog.dat
usr/local/bin/makegraphs
usr/local/bin/hddshutdown
usr/local/bin/vpn-watch
