# config/rootfiles/core/129/filelists/files
etc/system-release
etc/issue
var/ipfire/langs
etc/rc.d/init.d/firewall
etc/rc.d/init.d/network
etc/rc.d/init.d/networking/red.up/50-ipsec
srv/web/ipfire/cgi-bin/credits.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/netovpnsrv.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
usr/lib/firewall/firewall-lib.pl
usr/lib/firewall/ipsec-policy
usr/local/bin/ipsec-interfaces
usr/local/bin/ipsecctrl
