# config/rootfiles/core/158/filelists/files
etc/rc.d/init.d/firewall
srv/web/ipfire/cgi-bin/vpnmain.cgi
srv/web/ipfire/html/images/apple.png
