# config/rootfiles/core/158/filelists/files
etc/rc.d/init.d/firewall
etc/ssh/sshd_config
opt/pakfire/lib/functions.pl
opt/pakfire/pakfire
srv/web/ipfire/cgi-bin/dhcp.cgi
srv/web/ipfire/cgi-bin/memory.cgi
srv/web/ipfire/cgi-bin/pakfire.cgi
srv/web/ipfire/cgi-bin/traffic.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
srv/web/ipfire/html/images/apple.png
var/ipfire/backup/bin/backup.pl
var/ipfire/countries.pl
var/ipfire/general-functions.pl
var/ipfire/graphs.pl
var/ipfire/header.pl
var/ipfire/ids-functions.pl
var/ipfire/lang.pl
var/ipfire/location-functions.pl
var/ipfire/modem-lib.pl
