# config/rootfiles/core/124/filelists/files
etc/system-release
etc/issue
etc/rc.d/helper/aws-setup
etc/rc.d/init.d/aws
etc/rc.d/init.d/localnet
etc/rc.d/init.d/partresize
srv/web/ipfire/cgi-bin/firewall.cgi
srv/web/ipfire/cgi-bin/ids.cgi
usr/bin/install-bootloader
usr/local/bin/backupiso
var/ipfire/backup/exclude
var/ipfire/backup/include
var/ipfire/langs
