# config/rootfiles/core/124/filelists/files
etc/system-release
etc/issue
etc/nsswitch.conf
etc/rc.d/helper/aws-setup
etc/rc.d/init.d/apache
etc/rc.d/init.d/aws
etc/rc.d/init.d/localnet
etc/rc.d/init.d/networking/red.down/10-static-routes
etc/rc.d/init.d/partresize
etc/rc.d/init.d/static-routes
etc/ssh/ssh_config
etc/ssh/sshd_config
etc/sysctl.conf
etc/unbound/unbound.conf
opt/pakfire/lib/functions.pl
opt/pakfire/pakfire
srv/web/ipfire/cgi-bin/backup.cgi
srv/web/ipfire/cgi-bin/firewall.cgi
srv/web/ipfire/cgi-bin/fwhosts.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/pakfire.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
srv/web/ipfire/html/redirect-templates/legacy/template.html
usr/bin/install-bootloader
usr/local/bin/backupiso
usr/local/bin/rebuild-initrd
var/ipfire/backup/bin/backup.pl
var/ipfire/backup/exclude
var/ipfire/backup/include
var/ipfire/langs
