# config/rootfiles/core/124/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/localnet
var/ipfire/langs
