# config/rootfiles/core/1/files
etc/init.d/networking/red
var/ipfire/langs
srv/web/ipfire/cgi-bin/pakfire.cgi
srv/web/ipfire/cgi-bin/remote.cgi
usr/local/bin/sshctrl
usr/local/bin/redctrl
usr/local/bin/connscheduler
