# config/rootfiles/core/142/filelists/files
etc/system-release
etc/issue
srv/web/ipfire/cgi-bin/credits.cgi
var/ipfire/langs
etc/unbound/unbound.conf
etc/rc.d/helper/aws-setup
etc/rc.d/helper/azure-setup
etc/rc.d/init.d/suricata
etc/rc.d/init.d/unbound
etc/suricata/suricata.yaml
lib/udev/network-hotplug-bridges
opt/pakfire/etc/pakfire.conf
opt/pakfire/lib/functions.pl
srv/web/ipfire/cgi-bin/dns.cgi
srv/web/ipfire/cgi-bin/fireinfo.cgi
srv/web/ipfire/cgi-bin/pakfire.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
var/ipfire/backup/include
var/ipfire/suricata/ruleset-sources
