# config/rootfiles/core/196/filelists/files
etc/rc.d/init.d/networking/red
etc/rc.d/init.d/partresize
etc/rc.d/init.d/wireguard
srv/web/ipfire/cgi-bin/firewall.cgi
srv/web/ipfire/cgi-bin/fwhosts.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/netexternal.cgi
srv/web/ipfire/cgi-bin/netovpnsrv.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/pakfire.cgi
srv/web/ipfire/cgi-bin/vulnerabilities.cgi
srv/web/ipfire/cgi-bin/wireguard.cgi
usr/sbin/setup
var/ipfire/backup/bin/backup.pl
var/ipfire/graphs.pl
var/ipfire/header.pl
var/ipfire/ipblocklist-functions.pl
var/ipfire/langs/list
var/ipfire/wireguard-functions.pl
