# config/rootfiles/oldcore/77/filelists/files
etc/system-release
etc/issue
etc/collectd.conf
etc/rc.d/init.d/network
etc/sysctl.conf
run
srv/web/ipfire/cgi-bin/credits.cgi
srv/web/ipfire/cgi-bin/gui.cgi
srv/web/ipfire/cgi-bin/entropy.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/netinternal.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/upnp.cgi
srv/web/ipfire/cgi-bin/speed.cgi
srv/web/ipfire/cgi-bin/urlfilter.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
srv/web/ipfire/html/themes/darkdos
srv/web/ipfire/html/themes/ipfire
srv/web/ipfire/html/themes/ipfire-legacy
srv/web/ipfire/html/themes/maniac
usr/lib/firewall
usr/local/bin/backupiso
usr/local/bin/timezone-transition
usr/local/bin/setddns.pl
usr/sbin/updxlrator
var/ipfire/backup/bin/backup.pl
var/ipfire/backup/exclude
var/ipfire/backup/include
var/ipfire/general-functions.pl
var/ipfire/header.pl
var/ipfire/graphs.pl
var/ipfire/langs
var/ipfire/menu.d/10-system.menu
var/ipfire/menu.d/20-status.menu
