# config/rootfiles/oldcore/110/filelists/files
etc/system-release
etc/issue
etc/httpd/conf/server-tuning.conf
etc/rc.d/init.d/unbound
etc/sysctl.conf
opt/pakfire/lib/functions.pl
srv/web/ipfire/cgi-bin/entropy.cgi
srv/web/ipfire/cgi-bin/hardwaregraphs.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/netovpnsrv.cgi
srv/web/ipfire/cgi-bin/qos.cgi
srv/web/ipfire/cgi-bin/system.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
srv/web/ipfire/html/themes/darkdos/include/style.css
srv/web/ipfire/html/themes/ipfire/include/css/style.css
srv/web/ipfire/html/themes/ipfire-legacy/include/style.css
srv/web/ipfire/html/themes/maniac/include/style.css
usr/lib/firewall/ipsec-block
usr/lib/libssp.so.0
usr/lib/libssp.so.0.0.0
usr/local/bin/xt_geoip_update
usr/sbin/setup
var/ipfire/langs
var/ipfire/general-functions.pl
var/ipfire/graphs.pl
