# config/rootfiles/core/49/filelists/files
etc/system-release
etc/init.d/squid
srv/web/ipfire/cgi-bin/proxy.cgi
