# config/rootfiles/oldcore/117/filelists/files
etc/system-release
etc/issue
etc/httpd/conf/vhosts.d/ipfire-interface-ssl.conf
etc/rc.d/init.d/apache
etc/ssl/certs/ca-bundle.crt
etc/ssl/certs/ca-bundle.trust.crt
opt/pakfire/lib/functions.pl
srv/web/ipfire/cgi-bin/captive.cgi
srv/web/ipfire/cgi-bin/connections.cgi
srv/web/ipfire/cgi-bin/country.cgi
srv/web/ipfire/cgi-bin/credits.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/ipinfo.cgi
srv/web/ipfire/cgi-bin/logs.cgi/config.dat
srv/web/ipfire/cgi-bin/logs.cgi/firewalllog.dat
srv/web/ipfire/cgi-bin/logs.cgi/firewalllogcountry.dat
srv/web/ipfire/cgi-bin/logs.cgi/firewalllogip.dat
srv/web/ipfire/cgi-bin/logs.cgi/ids.dat
srv/web/ipfire/cgi-bin/logs.cgi/showrequestfromcountry.dat
srv/web/ipfire/cgi-bin/netexternal.cgi
srv/web/ipfire/cgi-bin/netother.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/routing.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
srv/web/ipfire/html/themes/darkdos/include/functions.pl
srv/web/ipfire/html/themes/ipfire-legacy/include/functions.pl
srv/web/ipfire/html/themes/ipfire/include/functions.pl
srv/web/ipfire/html/themes/maniac/include/functions.pl
usr/local/bin/captivectrl
usr/local/bin/syslogdctrl
var/ipfire/geoip-functions.pl
var/ipfire/network-functions.pl
