# config/rootfiles/core/127/filelists/files
etc/system-release
etc/issue
etc/modprobe.d/framebuffer.conf
srv/web/ipfire/cgi-bin/credits.cgi
var/ipfire/langs
etc/rc.d/init.d/aws
etc/rc.d/init.d/conntrackd
etc/rc.d/init.d/unbound
etc/rc.d/rc0.d/K77conntrackd
etc/rc.d/rc3.d/S22conntrackd
etc/rc.d/rc6.d/K77conntrackd
etc/sysctl.conf
srv/web/ipfire/cgi-bin/dnsforward.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
usr/lib/firewall/firewall-lib.pl
usr/local/bin/xt_geoip_update
var/ipfire/backup/include
var/ipfire/general-functions.pl
var/ipfire/geoip-functions.pl
