# config/rootfiles/core/kerneltest/filelists/files
etc/rc.d/init.d/firewall
etc/rc.d/init.d/leds
