# config/rootfiles/core/28/filelists/files
bin/probenic.sh
etc/ntp
etc/rc.d/init.d/beep
etc/rc.d/init.d/console
etc/rc.d/init.d/snort
etc/rc.d/init.d/networking/red
etc/ssh/moduli
etc/ssh/ssh_config
etc/ssh/sshd_config
etc/modprobe.d/blacklist
usr/bin/scp
usr/bin/sftp
usr/bin/slogin
usr/bin/ssh
usr/bin/ssh-add
usr/bin/ssh-agent
usr/bin/ssh-keygen
usr/bin/ssh-keyscan
usr/lib/openssh/sftp-server
usr/lib/openssh/ssh-keysign
usr/sbin/sshd
usr/lib/libcurl.so
usr/lib/libcurl.so.3
usr/lib/libcurl.so.4
usr/lib/libcurl.so.4.1.1
lib/modules/KVER-ipfire/misc/kqemu.ko
usr/local/bin/setaliases
srv/web/ipfire/cgi-bin/aliases.cgi
var/ipfire/menu.d/20-status.menu
var/ipfire/menu.d/30-network.menu
srv/web/ipfire/cgi-bin/hardwaregraphs.cgi
srv/web/ipfire/cgi-bin/credits.cgi
srv/web/ipfire/cgi-bin/backup.cgi
srv/web/ipfire/cgi-bin/ids.cgi
usr/local/bin/backupiso
