# config/rootfiles/core/168/filelists/files
etc/rc.d/init.d/networking/any
etc/rc.d/init.d/networking/red
srv/web/ipfire/cgi-bin/media.cgi
srv/web/ipfire/html/themes/ipfire/include/functions.pl
var/ipfire/header.pl
var/ipfire/menu.d/20-status.menu
var/ipfire/menu.d/30-network.menu
