# config/rootfiles/core/168/filelists/files
etc/rc.d/init.d/networking/any
etc/rc.d/init.d/networking/red
