# config/rootfiles/core/168/filelists/files
etc/rc.d/init.d/networking/any
etc/rc.d/init.d/networking/red
lib/firmware/amd/amd_sev_fam17h_model0xh.sbin
lib/firmware/amd/amd_sev_fam17h_model3xh.sbin
lib/firmware/amd/amd_sev_fam19h_model0xh.sbin
lib/firmware/amdgpu/aldebaran_mec2.bin
lib/firmware/amdgpu/aldebaran_mec.bin
lib/firmware/amdgpu/aldebaran_rlc.bin
lib/firmware/amdgpu/aldebaran_sdma.bin
lib/firmware/amdgpu/aldebaran_smc.bin
lib/firmware/amdgpu/aldebaran_sos.bin
lib/firmware/amdgpu/aldebaran_ta.bin
lib/firmware/amdgpu/aldebaran_vcn.bin
lib/firmware/amdgpu/arcturus_sdma.bin
lib/firmware/amdgpu/beige_goby_ce.bin
lib/firmware/amdgpu/beige_goby_dmcub.bin
lib/firmware/amdgpu/beige_goby_me.bin
lib/firmware/amdgpu/beige_goby_mec2.bin
lib/firmware/amdgpu/beige_goby_mec.bin
lib/firmware/amdgpu/beige_goby_pfp.bin
lib/firmware/amdgpu/beige_goby_rlc.bin
lib/firmware/amdgpu/beige_goby_sdma.bin
lib/firmware/amdgpu/beige_goby_smc.bin
lib/firmware/amdgpu/beige_goby_sos.bin
lib/firmware/amdgpu/beige_goby_ta.bin
lib/firmware/amdgpu/beige_goby_vcn.bin
lib/firmware/amdgpu/cyan_skillfish2_ce.bin
lib/firmware/amdgpu/cyan_skillfish2_me.bin
lib/firmware/amdgpu/cyan_skillfish2_mec2.bin
lib/firmware/amdgpu/cyan_skillfish2_mec.bin
lib/firmware/amdgpu/cyan_skillfish2_pfp.bin
lib/firmware/amdgpu/dcn_3_1_6_dmcub.bin
lib/firmware/amdgpu/dimgrey_cavefish_ce.bin
lib/firmware/amdgpu/dimgrey_cavefish_dmcub.bin
lib/firmware/amdgpu/dimgrey_cavefish_me.bin
lib/firmware/amdgpu/dimgrey_cavefish_mec2.bin
lib/firmware/amdgpu/dimgrey_cavefish_mec.bin
lib/firmware/amdgpu/dimgrey_cavefish_pfp.bin
lib/firmware/amdgpu/dimgrey_cavefish_rlc.bin
lib/firmware/amdgpu/dimgrey_cavefish_sdma.bin
lib/firmware/amdgpu/dimgrey_cavefish_smc.bin
lib/firmware/amdgpu/dimgrey_cavefish_sos.bin
lib/firmware/amdgpu/dimgrey_cavefish_ta.bin
lib/firmware/amdgpu/dimgrey_cavefish_vcn.bin
lib/firmware/amdgpu/gc_10_3_7_ce.bin
lib/firmware/amdgpu/gc_10_3_7_me.bin
lib/firmware/amdgpu/gc_10_3_7_mec2.bin
lib/firmware/amdgpu/gc_10_3_7_mec.bin
lib/firmware/amdgpu/gc_10_3_7_pfp.bin
lib/firmware/amdgpu/gc_10_3_7_rlc.bin
lib/firmware/amdgpu/green_sardine_asd.bin
lib/firmware/amdgpu/green_sardine_ce.bin
lib/firmware/amdgpu/green_sardine_dmcub.bin
lib/firmware/amdgpu/green_sardine_me.bin
lib/firmware/amdgpu/green_sardine_mec2.bin
lib/firmware/amdgpu/green_sardine_mec.bin
lib/firmware/amdgpu/green_sardine_pfp.bin
lib/firmware/amdgpu/green_sardine_rlc.bin
lib/firmware/amdgpu/green_sardine_ta.bin
lib/firmware/amdgpu/green_sardine_vcn.bin
lib/firmware/amdgpu/navi10_asd.bin
lib/firmware/amdgpu/navi10_ce.bin
lib/firmware/amdgpu/navi10_me.bin
lib/firmware/amdgpu/navi10_mec2.bin
lib/firmware/amdgpu/navi10_mec.bin
lib/firmware/amdgpu/navi10_pfp.bin
lib/firmware/amdgpu/navi10_sos.bin
lib/firmware/amdgpu/navi10_ta.bin
lib/firmware/amdgpu/navi10_vcn.bin
lib/firmware/amdgpu/navi12_asd.bin
lib/firmware/amdgpu/navi12_sdma1.bin
lib/firmware/amdgpu/navi12_sdma.bin
lib/firmware/amdgpu/navi12_sos.bin
lib/firmware/amdgpu/navi12_ta.bin
lib/firmware/amdgpu/navi12_vcn.bin
lib/firmware/amdgpu/navi14_asd.bin
lib/firmware/amdgpu/navi14_smc.bin
lib/firmware/amdgpu/navi14_sos.bin
lib/firmware/amdgpu/navi14_ta.bin
lib/firmware/amdgpu/navi14_vcn.bin
lib/firmware/amdgpu/navy_flounder_ce.bin
lib/firmware/amdgpu/navy_flounder_dmcub.bin
lib/firmware/amdgpu/navy_flounder_me.bin
lib/firmware/amdgpu/navy_flounder_mec2.bin
lib/firmware/amdgpu/navy_flounder_mec.bin
lib/firmware/amdgpu/navy_flounder_pfp.bin
lib/firmware/amdgpu/navy_flounder_rlc.bin
lib/firmware/amdgpu/navy_flounder_sdma.bin
lib/firmware/amdgpu/navy_flounder_smc.bin
lib/firmware/amdgpu/navy_flounder_sos.bin
lib/firmware/amdgpu/navy_flounder_ta.bin
lib/firmware/amdgpu/navy_flounder_vcn.bin
lib/firmware/amdgpu/picasso_asd.bin
lib/firmware/amdgpu/picasso_ce.bin
lib/firmware/amdgpu/picasso_me.bin
lib/firmware/amdgpu/picasso_mec2.bin
lib/firmware/amdgpu/picasso_mec.bin
lib/firmware/amdgpu/picasso_pfp.bin
lib/firmware/amdgpu/picasso_ta.bin
lib/firmware/amdgpu/picasso_vcn.bin
lib/firmware/amdgpu/psp_13_0_8_asd.bin
lib/firmware/amdgpu/psp_13_0_8_ta.bin
lib/firmware/amdgpu/psp_13_0_8_toc.bin
lib/firmware/amdgpu/raven2_asd.bin
lib/firmware/amdgpu/raven2_ce.bin
lib/firmware/amdgpu/raven2_me.bin
lib/firmware/amdgpu/raven2_mec2.bin
lib/firmware/amdgpu/raven2_mec.bin
lib/firmware/amdgpu/raven2_pfp.bin
lib/firmware/amdgpu/raven2_ta.bin
lib/firmware/amdgpu/raven2_vcn.bin
lib/firmware/amdgpu/raven_asd.bin
lib/firmware/amdgpu/raven_ce.bin
lib/firmware/amdgpu/raven_me.bin
lib/firmware/amdgpu/raven_mec2.bin
lib/firmware/amdgpu/raven_mec.bin
lib/firmware/amdgpu/raven_pfp.bin
lib/firmware/amdgpu/raven_ta.bin
lib/firmware/amdgpu/raven_vcn.bin
lib/firmware/amdgpu/renoir_asd.bin
lib/firmware/amdgpu/renoir_ce.bin
lib/firmware/amdgpu/renoir_dmcub.bin
lib/firmware/amdgpu/renoir_me.bin
lib/firmware/amdgpu/renoir_mec2.bin
lib/firmware/amdgpu/renoir_mec.bin
lib/firmware/amdgpu/renoir_pfp.bin
lib/firmware/amdgpu/renoir_rlc.bin
lib/firmware/amdgpu/renoir_ta.bin
lib/firmware/amdgpu/renoir_vcn.bin
lib/firmware/amdgpu/sdma_5_2_7.bin
lib/firmware/amdgpu/sienna_cichlid_ce.bin
lib/firmware/amdgpu/sienna_cichlid_dmcub.bin
lib/firmware/amdgpu/sienna_cichlid_me.bin
lib/firmware/amdgpu/sienna_cichlid_mec2.bin
lib/firmware/amdgpu/sienna_cichlid_mec.bin
lib/firmware/amdgpu/sienna_cichlid_pfp.bin
lib/firmware/amdgpu/sienna_cichlid_rlc.bin
lib/firmware/amdgpu/sienna_cichlid_sdma.bin
lib/firmware/amdgpu/sienna_cichlid_smc.bin
lib/firmware/amdgpu/sienna_cichlid_sos.bin
lib/firmware/amdgpu/sienna_cichlid_ta.bin
lib/firmware/amdgpu/sienna_cichlid_vcn.bin
lib/firmware/amdgpu/vangogh_asd.bin
lib/firmware/amdgpu/vangogh_ce.bin
lib/firmware/amdgpu/vangogh_me.bin
lib/firmware/amdgpu/vangogh_mec2.bin
lib/firmware/amdgpu/vangogh_mec.bin
lib/firmware/amdgpu/vangogh_pfp.bin
lib/firmware/amdgpu/vega10_asd.bin
lib/firmware/amdgpu/vega10_ce.bin
lib/firmware/amdgpu/vega10_me.bin
lib/firmware/amdgpu/vega10_mec2.bin
lib/firmware/amdgpu/vega10_mec.bin
lib/firmware/amdgpu/vega10_pfp.bin
lib/firmware/amdgpu/vega12_asd.bin
lib/firmware/amdgpu/vega12_ce.bin
lib/firmware/amdgpu/vega12_me.bin
lib/firmware/amdgpu/vega12_mec2.bin
lib/firmware/amdgpu/vega12_mec.bin
lib/firmware/amdgpu/vega12_pfp.bin
lib/firmware/amdgpu/vega20_asd.bin
lib/firmware/amdgpu/vega20_ce.bin
lib/firmware/amdgpu/vega20_me.bin
lib/firmware/amdgpu/vega20_mec2.bin
lib/firmware/amdgpu/vega20_mec.bin
lib/firmware/amdgpu/vega20_pfp.bin
lib/firmware/amdgpu/yellow_carp_asd.bin
lib/firmware/amdgpu/yellow_carp_ce.bin
lib/firmware/amdgpu/yellow_carp_dmcub.bin
lib/firmware/amdgpu/yellow_carp_me.bin
lib/firmware/amdgpu/yellow_carp_mec2.bin
lib/firmware/amdgpu/yellow_carp_mec.bin
lib/firmware/amdgpu/yellow_carp_pfp.bin
lib/firmware/amdgpu/yellow_carp_rlc.bin
lib/firmware/amdgpu/yellow_carp_ta.bin
lib/firmware/amd-ucode/microcode_amd_fam17h.bin
lib/firmware/amd-ucode/microcode_amd_fam19h.bin
lib/firmware/amphion/vpu/vpu_fw_imx8_dec.bin
lib/firmware/amphion/vpu/vpu_fw_imx8_enc.bin
lib/firmware/ath11k/WCN6855/hw2.0/amss.bin
lib/firmware/ath11k/WCN6855/hw2.0/board-2.bin
lib/firmware/ath11k/WCN6855/hw2.0/m3.bin
lib/firmware/ath11k/WCN6855/hw2.0/regdb.bin
lib/firmware/ath11k/WCN6855/hw2.1/amss.bin
lib/firmware/ath11k/WCN6855/hw2.1/board-2.bin
lib/firmware/ath11k/WCN6855/hw2.1/m3.bin
lib/firmware/ath11k/WCN6855/hw2.1/regdb.bin
lib/firmware/brcm/brcmfmac43430-sdio.friendlyarm,nanopi-r1.txt
lib/firmware/brcm/brcmfmac4356-sdio.AP6356S.txt
lib/firmware/brcm/brcmfmac4356-sdio.firefly,firefly-rk3399.txt
lib/firmware/cnm/wave521c_j721s2_codec_fw.bin
lib/firmware/cxgb4/t4fw-1.26.6.0.bin
lib/firmware/cxgb4/t5fw-1.26.6.0.bin
lib/firmware/cxgb4/t6fw-1.26.6.0.bin
lib/firmware/i915/adlp_dmc_ver2_16.bin
lib/firmware/i915/adlp_guc_69.0.3.bin
lib/firmware/i915/bxt_guc_69.0.3.bin
lib/firmware/i915/cml_guc_69.0.3.bin
lib/firmware/i915/dg1_guc_69.0.3.bin
lib/firmware/i915/ehl_guc_69.0.3.bin
lib/firmware/i915/glk_guc_69.0.3.bin
lib/firmware/i915/icl_guc_69.0.3.bin
lib/firmware/i915/kbl_guc_69.0.3.bin
lib/firmware/i915/skl_guc_69.0.3.bin
lib/firmware/i915/tgl_guc_69.0.3.bin
lib/firmware/intel/ibt-0040-0041.sfi
lib/firmware/intel/ibt-0040-1020.sfi
lib/firmware/intel/ibt-0040-2120.sfi
lib/firmware/intel/ibt-0040-4150.sfi
lib/firmware/intel/ibt-0041-0041.sfi
lib/firmware/intel/ibt-1040-0041.sfi
lib/firmware/intel/ibt-1040-1020.sfi
lib/firmware/intel/ibt-1040-2120.sfi
lib/firmware/intel/ibt-1040-4150.sfi
lib/firmware/intel/ibt-11-5.sfi
lib/firmware/intel/ibt-12-16.sfi
lib/firmware/intel/ibt-17-16-1.sfi
lib/firmware/intel/ibt-17-2.sfi
lib/firmware/intel/ibt-18-16-1.sfi
lib/firmware/intel/ibt-18-2.sfi
lib/firmware/intel/ibt-19-0-0.sfi
lib/firmware/intel/ibt-19-0-1.sfi
lib/firmware/intel/ibt-19-0-4.sfi
lib/firmware/intel/ibt-19-16-4.sfi
lib/firmware/intel/ibt-19-240-1.sfi
lib/firmware/intel/ibt-19-240-4.sfi
lib/firmware/intel/ibt-19-32-0.sfi
lib/firmware/intel/ibt-19-32-1.sfi
lib/firmware/intel/ibt-19-32-4.sfi
lib/firmware/intel/ibt-20-0-3.sfi
lib/firmware/intel/ibt-20-1-3.sfi
lib/firmware/intel/ibt-20-1-4.sfi
lib/firmware/intel/ibt-hw-37.8.10-fw-22.50.19.14.f.bseq
lib/firmware/intel/ice/ddp-comms/ice_comms-1.3.31.0.pkg
lib/firmware/intel/ice/ddp-wireless_edge
lib/firmware/intel/ice/ddp-wireless_edge/ice_wireless_edge-1.3.7.0.pkg
lib/firmware/iwlwifi-9000-pu-b0-jf-b0-46.ucode
lib/firmware/iwlwifi-9260-th-b0-jf-b0-46.ucode
lib/firmware/iwlwifi-cc-a0-66.ucode
lib/firmware/iwlwifi-cc-a0-68.ucode
lib/firmware/iwlwifi-cc-a0-71.ucode
lib/firmware/iwlwifi-Qu-b0-hr-b0-66.ucode
lib/firmware/iwlwifi-Qu-b0-hr-b0-68.ucode
lib/firmware/iwlwifi-Qu-b0-hr-b0-71.ucode
lib/firmware/iwlwifi-Qu-b0-jf-b0-66.ucode
lib/firmware/iwlwifi-Qu-b0-jf-b0-68.ucode
lib/firmware/iwlwifi-Qu-b0-jf-b0-71.ucode
lib/firmware/iwlwifi-Qu-c0-hr-b0-66.ucode
lib/firmware/iwlwifi-Qu-c0-hr-b0-68.ucode
lib/firmware/iwlwifi-Qu-c0-hr-b0-71.ucode
lib/firmware/iwlwifi-Qu-c0-jf-b0-66.ucode
lib/firmware/iwlwifi-Qu-c0-jf-b0-68.ucode
lib/firmware/iwlwifi-Qu-c0-jf-b0-71.ucode
lib/firmware/iwlwifi-QuZ-a0-hr-b0-66.ucode
lib/firmware/iwlwifi-QuZ-a0-hr-b0-68.ucode
lib/firmware/iwlwifi-QuZ-a0-hr-b0-71.ucode
lib/firmware/iwlwifi-QuZ-a0-jf-b0-66.ucode
lib/firmware/iwlwifi-QuZ-a0-jf-b0-68.ucode
lib/firmware/iwlwifi-QuZ-a0-jf-b0-71.ucode
lib/firmware/iwlwifi-so-a0-gf4-a0-68.ucode
lib/firmware/iwlwifi-so-a0-gf4-a0-71.ucode
lib/firmware/iwlwifi-so-a0-gf4-a0.pnvm
lib/firmware/iwlwifi-so-a0-gf-a0-68.ucode
lib/firmware/iwlwifi-so-a0-gf-a0-71.ucode
lib/firmware/iwlwifi-so-a0-gf-a0.pnvm
lib/firmware/iwlwifi-so-a0-hr-b0-68.ucode
lib/firmware/iwlwifi-so-a0-hr-b0-71.ucode
lib/firmware/iwlwifi-so-a0-jf-b0-68.ucode
lib/firmware/iwlwifi-so-a0-jf-b0-71.ucode
lib/firmware/iwlwifi-ty-a0-gf-a0-66.ucode
lib/firmware/iwlwifi-ty-a0-gf-a0-68.ucode
lib/firmware/iwlwifi-ty-a0-gf-a0-71.ucode
lib/firmware/iwlwifi-ty-a0-gf-a0.pnvm
lib/firmware/mediatek/BT_RAM_CODE_MT7961_1_2_hdr.bin
lib/firmware/mediatek/mt7915_eeprom.bin
lib/firmware/mediatek/mt7915_eeprom_dbdc.bin
lib/firmware/mediatek/mt7915_rom_patch.bin
lib/firmware/mediatek/mt7915_wa.bin
lib/firmware/mediatek/mt7915_wm.bin
lib/firmware/mediatek/mt7916_eeprom.bin
lib/firmware/mediatek/mt7916_rom_patch.bin
lib/firmware/mediatek/mt7916_wa.bin
lib/firmware/mediatek/mt7916_wm.bin
lib/firmware/mediatek/mt7986_eeprom_mt7975_dual.bin
lib/firmware/mediatek/mt7986_eeprom_mt7976.bin
lib/firmware/mediatek/mt7986_eeprom_mt7976_dbdc.bin
lib/firmware/mediatek/mt7986_eeprom_mt7976_dual.bin
lib/firmware/mediatek/mt7986_rom_patch.bin
lib/firmware/mediatek/mt7986_rom_patch_mt7975.bin
lib/firmware/mediatek/mt7986_wa.bin
lib/firmware/mediatek/mt7986_wm.bin
lib/firmware/mediatek/mt7986_wm_mt7975.bin
lib/firmware/mediatek/mt8173/vpu_d.bin
lib/firmware/mediatek/mt8173/vpu_p.bin
lib/firmware/mediatek/WIFI_MT7961_patch_mcu_1_2_hdr.bin
lib/firmware/mediatek/WIFI_RAM_CODE_MT7961_1.bin
lib/firmware/mrvl/cpt01/ae.out
lib/firmware/mrvl/cpt01/ie.out
lib/firmware/mrvl/cpt01/se.out
lib/firmware/mrvl/cpt02/ae.out
lib/firmware/mrvl/cpt02/ie.out
lib/firmware/mrvl/cpt02/se.out
lib/firmware/mrvl/cpt03/ae.out
lib/firmware/mrvl/cpt03/ie.out
lib/firmware/mrvl/cpt03/se.out
lib/firmware/mrvl/cpt04/ae.out
lib/firmware/mrvl/cpt04/ie.out
lib/firmware/mrvl/cpt04/se.out
lib/firmware/nvidia/ga102/acr/ucode_ahesasc.bin
lib/firmware/nvidia/ga102/acr/ucode_asb.bin
lib/firmware/nvidia/ga102/acr/ucode_unload.bin
lib/firmware/nvidia/ga102/gr/fecs_bl.bin
lib/firmware/nvidia/ga102/gr/fecs_sig.bin
lib/firmware/nvidia/ga102/gr/gpccs_bl.bin
lib/firmware/nvidia/ga102/gr/gpccs_sig.bin
lib/firmware/nvidia/ga102/gr/NET_img.bin
lib/firmware/nvidia/ga102/nvdec/scrubber.bin
lib/firmware/nvidia/ga102/sec2/desc.bin
lib/firmware/nvidia/ga102/sec2/hs_bl_sig.bin
lib/firmware/nvidia/ga102/sec2/image.bin
lib/firmware/nvidia/ga102/sec2/sig.bin
lib/firmware/nvidia/ga103/acr/ucode_ahesasc.bin
lib/firmware/nvidia/ga103/acr/ucode_asb.bin
lib/firmware/nvidia/ga103/acr/ucode_unload.bin
lib/firmware/nvidia/ga103/gr/fecs_bl.bin
lib/firmware/nvidia/ga103/gr/fecs_sig.bin
lib/firmware/nvidia/ga103/gr/gpccs_bl.bin
lib/firmware/nvidia/ga103/gr/gpccs_sig.bin
lib/firmware/nvidia/ga103/gr/NET_img.bin
lib/firmware/nvidia/ga103/nvdec/scrubber.bin
lib/firmware/nvidia/ga103/sec2/desc.bin
lib/firmware/nvidia/ga103/sec2/hs_bl_sig.bin
lib/firmware/nvidia/ga103/sec2/image.bin
lib/firmware/nvidia/ga103/sec2/sig.bin
lib/firmware/nvidia/ga104/acr/ucode_ahesasc.bin
lib/firmware/nvidia/ga104/acr/ucode_asb.bin
lib/firmware/nvidia/ga104/acr/ucode_unload.bin
lib/firmware/nvidia/ga104/gr/fecs_bl.bin
lib/firmware/nvidia/ga104/gr/fecs_sig.bin
lib/firmware/nvidia/ga104/gr/gpccs_bl.bin
lib/firmware/nvidia/ga104/gr/gpccs_sig.bin
lib/firmware/nvidia/ga104/gr/NET_img.bin
lib/firmware/nvidia/ga104/nvdec/scrubber.bin
lib/firmware/nvidia/ga104/sec2/desc.bin
lib/firmware/nvidia/ga104/sec2/hs_bl_sig.bin
lib/firmware/nvidia/ga104/sec2/image.bin
lib/firmware/nvidia/ga104/sec2/sig.bin
lib/firmware/nvidia/ga106/acr/ucode_ahesasc.bin
lib/firmware/nvidia/ga106/acr/ucode_asb.bin
lib/firmware/nvidia/ga106/acr/ucode_unload.bin
lib/firmware/nvidia/ga106/gr/fecs_bl.bin
lib/firmware/nvidia/ga106/gr/fecs_sig.bin
lib/firmware/nvidia/ga106/gr/gpccs_bl.bin
lib/firmware/nvidia/ga106/gr/gpccs_sig.bin
lib/firmware/nvidia/ga106/gr/NET_img.bin
lib/firmware/nvidia/ga106/nvdec/scrubber.bin
lib/firmware/nvidia/ga106/sec2/desc.bin
lib/firmware/nvidia/ga106/sec2/hs_bl_sig.bin
lib/firmware/nvidia/ga106/sec2/image.bin
lib/firmware/nvidia/ga106/sec2/sig.bin
lib/firmware/nvidia/ga107/acr/ucode_ahesasc.bin
lib/firmware/nvidia/ga107/acr/ucode_asb.bin
lib/firmware/nvidia/ga107/acr/ucode_unload.bin
lib/firmware/nvidia/ga107/gr/fecs_bl.bin
lib/firmware/nvidia/ga107/gr/fecs_sig.bin
lib/firmware/nvidia/ga107/gr/gpccs_bl.bin
lib/firmware/nvidia/ga107/gr/gpccs_sig.bin
lib/firmware/nvidia/ga107/gr/NET_img.bin
lib/firmware/nvidia/ga107/nvdec/scrubber.bin
lib/firmware/nvidia/ga107/sec2/desc.bin
lib/firmware/nvidia/ga107/sec2/hs_bl_sig.bin
lib/firmware/nvidia/ga107/sec2/image.bin
lib/firmware/nvidia/ga107/sec2/sig.bin
lib/firmware/qca/nvm_usb_00130201_010a.bin
lib/firmware/qca/nvm_usb_00130201_010b.bin
lib/firmware/qca/nvm_usb_00130201_0303.bin
lib/firmware/qca/nvm_usb_00130201.bin
lib/firmware/qca/nvm_usb_00130201_gf_010a.bin
lib/firmware/qca/nvm_usb_00130201_gf_010b.bin
lib/firmware/qca/nvm_usb_00130201_gf_0303.bin
lib/firmware/qca/nvm_usb_00130201_gf.bin
lib/firmware/qca/rampatch_usb_00130200.bin
lib/firmware/qca/rampatch_usb_00130201.bin
lib/firmware/rtl_bt/rtl8822cs_fw.bin
lib/firmware/rtl_bt/rtl8822cu_fw.bin
lib/firmware/rtl_bt/rtl8852au_fw.bin
lib/firmware/rtl_bt/rtl8852bu_config.bin
lib/firmware/rtl_bt/rtl8852bu_fw.bin
lib/firmware/rtw88/rtw8821c_fw.bin
lib/firmware/rtw88/rtw8822c_fw.bin
lib/firmware/rtw89/rtw8852a_fw.bin
lib/firmware/wfx/wfm_wf200_C0.sec
usr/lib/firewall/rules.pl
srv/web/ipfire/cgi-bin/media.cgi
srv/web/ipfire/cgi-bin/wakeonlan.cgi
srv/web/ipfire/html/themes/ipfire/include/functions.pl
var/ipfire/header.pl
var/ipfire/menu.d/20-status.menu
var/ipfire/menu.d/30-network.menu
