# config/rootfiles/core/137/filelists/files
etc/system-release
etc/issue
srv/web/ipfire/cgi-bin/credits.cgi
var/ipfire/langs
etc/logrotate.conf
etc/rc.d/init.d/firewall
etc/rc.d/init.d/unbound
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
usr/lib/firewall/rules.pl
usr/sbin/firewall-policy
var/ipfire/backup/bin/backup.pl
var/ipfire/suricata/ruleset-sources
