# config/rootfiles/oldcore/113/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/ntp
etc/rc.d/init.d/unbound
var/ipfire/langs
