# config/rootfiles/core/140/filelists/files
etc/system-release
etc/issue
srv/web/ipfire/cgi-bin/credits.cgi
var/ipfire/langs
etc/rc.d/init.d/networking/red.up/06-safe-search
etc/rc.d/init.d/unbound
etc/rc.d/init.d/suricata
opt/pakfire/lib/functions.pl
srv/web/ipfire/cgi-bin/ids.cgi
