# config/rootfiles/core/140/filelists/files
etc/system-release
etc/issue
srv/web/ipfire/cgi-bin/credits.cgi
var/ipfire/langs
etc/rc.d/helper/aws-setup
etc/rc.d/helper/azure-setup
etc/rc.d/init.d/lvmetad
etc/rc.d/init.d/networking/dhcpcd.exe
etc/rc.d/init.d/networking/red
etc/rc.d/init.d/networking/red.up/05-update-dns-forwarders
etc/rc.d/init.d/partresize
etc/rc.d/init.d/suricata
etc/rc.d/init.d/unbound
etc/rc.d/rcsysinit.d/S09lvmetad
etc/ppp/ip-down
etc/ppp/ip-up
etc/unbound/unbound.conf
opt/pakfire/etc/pakfire.conf
opt/pakfire/lib/functions.pl
srv/web/ipfire/cgi-bin/dns.cgi
srv/web/ipfire/cgi-bin/dnsforward.cgi
srv/web/ipfire/cgi-bin/gui.cgi
srv/web/ipfire/cgi-bin/hosts.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/netexternal.cgi
srv/web/ipfire/cgi-bin/pppsetup.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
usr/local/bin/backupiso
usr/local/bin/convert-dns-settings
usr/local/bin/unboundctrl
usr/sbin/convert-snort
usr/sbin/setup
usr/sbin/unbound-dhcp-leases-bridge
var/ipfire/backup/bin/backup.pl
var/ipfire/ids-functions.pl
var/ipfire/menu.d/30-network.menu
var/ipfire/general-functions.pl
