# config/rootfiles/oldcore/26/filelists/files
usr/sbin/pptp
etc/rc.d/init.d/networking/red
srv/web/ipfire/cgi-bin/pppsetup.cgi
srv/web/ipfire/html/include/jquery-1.2.6.min.js
srv/web/ipfire/html/themes/ipfire/include/functions.pl
srv/web/ipfire/cgi-bin/proxy.cgi
usr/sbin/updxlrator
usr/lib/apache/mod_proxy.so
usr/lib/apache/mod_proxy_ajp.so
usr/lib/apache/mod_proxy_balancer.so
usr/lib/apache/mod_proxy_connect.so
usr/lib/apache/mod_proxy_ftp.so
usr/lib/apache/mod_proxy_http.so
etc/ssl/openssl.cnf
usr/bin/openssl
usr/lib/libcrypto.so
usr/lib/libcrypto.so.0.9.8
usr/lib/libssl.so
usr/lib/libssl.so.0.9.8
usr/bin/squidGuard
var/lib/squidguard
usr/sbin/redirect_wrapper
