# config/rootfiles/updater/filelists/files
opt/pakfire/db/core/mine
etc/logrotate.conf
etc/rc.d/init.d/mISDN
etc/modprobe.d/blacklist
usr/bin/bc
usr/bin/dc
etc/ppp/dialer
etc/rc.d/init.d/checkfstab
etc/rc.d/rcsysinit.d/S19checkfstab
etc/rc.d/init.d/dhcp
etc/rc.d/init.d/tmpfs
etc/rc.d/rc0.d/K85tmpfs
etc/rc.d/rc3.d/S01tmpfs
etc/rc.d/rc6.d/K85tmpfs
etc/rc.d/init.d/dnsmasq
etc/rc.d/init.d/sysklogd
etc/rc.d/init.d/network
etc/rc.d/init.d/networking/any
etc/rc.d/init.d/networking/red
etc/rc.d/init.d/networking/red.up/40-ipac
etc/udev/dvb.sh
etc/udev/rules.d/10-dvb.rules
etc/udev/rules.d/55-scsi-cdrom.rules
etc/linuxigd/gatedesc.xml
usr/local/bin/getiptstate
usr/local/bin/makegraphs
usr/local/bin/hddshutdown
usr/local/bin/squidctrl
usr/local/bin/setddns.pl
var/ipfire/backup/bin/backup.pl
var/ipfire/sensors/
var/ipfire/menu.d/20-status.menu
var/ipfire/graphs.pl
var/ipfire/firebuild
var/ipfire/qos/bin/makeqosscripts.pl
var/ipfire/proxy/advanced/useragents
srv/web/ipfire/html/images/updbooster/
var/ipfire/updatexlrator/bin/
usr/sbin/updxlrator
var/ipfire/langs
usr/local/bin/qosd
usr/local/bin/backupiso
usr/bin/mkisofs
usr/bin/cdrecord
