# config/rootfiles/updater/filelists/files
opt/pakfire/db/core/mine
etc/init.d/mISDN
etc/modprobe.d/blacklist
usr/bin/bc
usr/bin/dc
etc/fstab
etc/ppp/dialer
etc/rc.d/init.d/checkfstab
etc/rc.d/rcsysinit.d/S19checkfstab
etc/rc.d/init.d/dhcp
etc/rc.d/init.d/dnsmasq
etc/rc.d/init.d/network
etc/rc.d/init.d/networking/any
etc/rc.d/init.d/networking/red
etc/rc.d/init.d/networking/red.up/40-ipac
etc/udev/dvb.sh
etc/udev/rules.d/10-dvb.rules
etc/udev/rules.d/55-scsi-cdrom.rules
etc/linuxigd/gatedesc.xml
usr/local/bin/getiptstate
usr/local/bin/makegraphs
usr/local/bin/squidctrl
usr/local/bin/setddns.pl
srv/web/ipfire/cgi-bin/credits.cgi
srv/web/ipfire/cgi-bin/speed.cgi
srv/web/ipfire/cgi-bin/traffic.cgi
srv/web/ipfire/cgi-bin/hardwaregraphs.cgi
srv/web/ipfire/cgi-bin/connections.cgi
srv/web/ipfire/cgi-bin/graphs.cgi
srv/web/ipfire/cgi-bin/pppsetup.cgi
srv/web/ipfire/cgi-bin/netexternal.cgi
srv/web/ipfire/cgi-bin/netinternal.cgi
srv/web/ipfire/cgi-bin/netother.cgi
srv/web/ipfire/cgi-bin/system.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/ddns.cgi
var/ipfire/sensors/
var/ipfire/graphs.pl
var/ipfire/firebuild
var/ipfire/proxy/advanced/useragents
