# config/rootfiles/core/9/files
lib/modules/2.6.16.57-ipfire/misc/kqemu.ko
lib/modules/2.6.16.57-ipfire-smp/misc/kqemu.ko
sbin/hdparm
usr/lib/collectd
usr/sbin/collectd
usr/sbin/collectdmon
etc/init.d/networking/red
etc/squid
sbin/capiinit
usr/bin/capiinfo
usr/lib/libcapi20.so
usr/lib/libcapi20.so.3
usr/lib/libcapi20.so.3.0.4
usr/lib/squid
usr/sbin/squid
var/ipfire/graphs.pl
usr/local/bin/addonctrl
srv/web/ipfire/cgi-bin/services.cgi
srv/web/ipfire/cgi-bin/system.cgi
var/ipfire/menu.d/20-status.menu
var/ipfire/langs
