# config/rootfiles/core/47/filelists/files
etc/system-release
etc/rc.d/init.d/sshd
srv/web/ipfire/cgi-bin/services.cgi
srv/web/ipfire/cgi-bin/remote.cgi
