# config/rootfiles/core/176/filelists/files
etc/rc.d/init.d/partresize
lib/udev/rules.d/69-dm-lvm.rules
srv/web/ipfire/cgi-bin/extrahd.cgi
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
var/ipfire/qos/bin/parse-func.pl
