# config/rootfiles/core/176/filelists/files
etc/rc.d/init.d/partresize
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
