# config/rootfiles/core/27/files
etc/sysctl.conf
etc/squid/errors
usr/bin/squidclient
usr/lib/squid/auth
usr/lib/squid/cachemgr.cgi
usr/lib/squid/diskd-daemon
usr/lib/squid/errors/English
usr/lib/squid/errors/German
usr/lib/squid/icons
usr/lib/squid/logfile-daemon
usr/lib/squid/mib.txt
usr/lib/squid/msnt_auth
usr/lib/squid/ncsa_auth
usr/lib/squid/ntlm_auth
usr/lib/squid/smb_auth
usr/lib/squid/smb_auth.pl
usr/lib/squid/smb_auth.sh
usr/lib/squid/squid_ldap_auth
usr/lib/squid/unlinkd
usr/sbin/squid
var/ipfire/updatexlrator/bin/checkdeaddl
var/ipfire/updatexlrator/bin/checkup
var/ipfire/updatexlrator/bin/convert
var/ipfire/updatexlrator/bin/download
var/ipfire/updatexlrator/bin/lscache
var/ipfire/updatexlrator/updxlrator-lib.pl
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/updatexlrator.cgi
var/ipfire/langs
lib/modules/2.6.25.19-ipfire/kernel/drivers/net/imq.ko
var/ipfire/qos/bin/makeqosscripts.pl
etc/init.d/network
