# config/rootfiles/oldcore/131/filelists/files
etc/system-release
etc/issue
srv/web/ipfire/cgi-bin/credits.cgi
etc/collectd.conf
etc/logrotate.conf
etc/rc.d/init.d/collectd
etc/rc.d/init.d/firewall
etc/rc.d/init.d/networking/red.up/23-suricata
etc/rc.d/init.d/suricata
etc/rc.d/init.d/unbound
etc/syslog.conf
lib/udev/network-hotplug-rename
opt/pakfire/etc/pakfire.conf
srv/web/ipfire/cgi-bin/aliases.cgi
srv/web/ipfire/cgi-bin/dnsforward.cgi
srv/web/ipfire/cgi-bin/hosts.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/logs.cgi/ids.dat
srv/web/ipfire/cgi-bin/logs.cgi/log.dat
srv/web/ipfire/cgi-bin/ovpnmain.cgi
srv/web/ipfire/cgi-bin/remote.cgi
srv/web/ipfire/cgi-bin/services.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
usr/local/bin/ipsec-interfaces
usr/local/bin/sshctrl
usr/local/bin/suricatactrl
usr/local/bin/update-ids-ruleset
usr/sbin/convert-snort
usr/sbin/unbound-dhcp-leases-bridge
usr/sbin/setup
var/ipfire/backup/bin/backup.pl
var/ipfire/backup/include
var/ipfire/general-functions.pl
var/ipfire/geoip-functions.pl
var/ipfire/ids-functions.pl
var/ipfire/langs
var/ipfire/menu.d/40-services.menu
var/ipfire/menu.d/50-firewall.menu
