# config/rootfiles/oldcore/189/filelists/files
etc/rc.d/init.d/functions
var/ipfire/header.pl
