# config/rootfiles/oldcore/57/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/networking/red.up/10-static-routes
etc/rc.d/init.d/networking/red.down/10-static-routes
etc/rc.d/init.d/networking/any
