# config/rootfiles/oldcore/178/filelists/files
etc/rc.d/init.d/smt
srv/web/ipfire/cgi-bin/vulnerabilities.cgi
