# config/rootfiles/core/52/filelists/files
etc/system-release
var/ipfire/langs/
etc/rc.d/init.d/console
usr/local/sbin/setup
usr/local/bin/setddns.pl
srv/web/ipfire/cgi-bin/ddns.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
var/ipfire/graphs.pl
