# config/rootfiles/core/52/filelists/files
etc/system-release
var/ipfire/langs/
etc/rc.d/init.d/console
usr/local/sbin/setup
