# config/rootfiles/core/116/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/apache
srv/web/ipfire/cgi-bin/proxy.cgi
