# config/rootfiles/core/44/filelists/files
etc/system-release
etc/rc.d/init.d/checkfs
etc/rc.d/init.d/checkfstab
etc/rc.d/init.d/firewall
etc/rc.d/init.d/squid
etc/rc.d/init.d/networking/dhcpcd.exe
srv/web/ipfire/cgi-bin/extrahd.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/pakfire.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/speed.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/html/images/urlfilter/copy.gif
var/ipfire/general-functions.pl
var/ipfire/graphs.pl
var/ipfire/extrahd/bin/extrahd.pl
var/ipfire/langs/de.pl
var/ipfire/langs/en.pl
var/ipfire/langs/fr.pl
var/ipfire/langs/es.pl
usr/local/bin/timectrl
usr/local/bin/rebuild-initrd
usr/local/bin/scanhd
usr/local/sbin/setup
usr/bin/ldd
usr/share/doc/licenses/GPLv3
