# config/rootfiles/core/44/filelists/files
etc/system-release
etc/issue
etc/rc.d/init.d/checkfs
etc/rc.d/init.d/checkfstab
etc/rc.d/init.d/fireinfo
etc/rc.d/init.d/firewall
etc/rc.d/init.d/squid
etc/rc.d/init.d/networking/dhcpcd.exe
etc/rc.d/init.d/networking/red
etc/rc.d/init.d/networking/red.up/99-fireinfo
etc/rc.d/rc3.d/S15fireinfo
srv/web/ipfire/cgi-bin/fireinfo.cgi
srv/web/ipfire/cgi-bin/extrahd.cgi
srv/web/ipfire/cgi-bin/credits.cgi
srv/web/ipfire/cgi-bin/index.cgi
srv/web/ipfire/cgi-bin/ipinfo.cgi
srv/web/ipfire/cgi-bin/pakfire.cgi
srv/web/ipfire/cgi-bin/proxy.cgi
srv/web/ipfire/cgi-bin/speed.cgi
srv/web/ipfire/cgi-bin/vpnmain.cgi
srv/web/ipfire/cgi-bin/ids.cgi
srv/web/ipfire/cgi-bin/outgoingfw.cgi
srv/web/ipfire/cgi-bin/outgoinggrp.cgi
srv/web/ipfire/cgi-bin/pppsetup.cgi
srv/web/ipfire/cgi-bin/hardwaregraphs.cgi
srv/web/ipfire/cgi-bin/memory.cgi
srv/web/ipfire/cgi-bin/netexternal.cgi
srv/web/ipfire/cgi-bin/netinternal.cgi
srv/web/ipfire/cgi-bin/system.cgi
srv/web/ipfire/cgi-bin/services.cgi
srv/web/ipfire/cgi-bin/backup.cgi
srv/web/ipfire/cgi-bin/logs.cgi/log.dat
srv/web/ipfire/html/images/urlfilter/copy.gif
var/ipfire/general-functions.pl
var/ipfire/graphs.pl
var/ipfire/extrahd/bin/extrahd.pl
var/ipfire/langs/de.pl
var/ipfire/langs/en.pl
var/ipfire/langs/fr.pl
var/ipfire/langs/es.pl
var/ipfire/menu.d/10-system.menu
usr/local/bin/timectrl
usr/local/bin/rebuild-initrd
usr/local/bin/scanhd
usr/local/sbin/setup
usr/local/share/GeoIP/GeoIP.dat
usr/bin/ldd
usr/share/doc/licenses/GPLv3
